# objectives/persistence/callback/polling/javascript.yaml
# Persistent Polling - JavaScript/Node.js
# Detects setInterval/setTimeout patterns used for persistent callbacks
# ATT&CK: T1053 (Scheduled Task/Job)

traits:
  # === setInterval HTTP atomic indicators ===
  - id: setInterval-http-call
    desc: setInterval with http
    crit: notable
    conf: 0.6
    for: [javascript]
    if:
      type: string
      regex: "setInterval.{0,50}http"

  - id: setInterval-axios
    desc: setInterval with axios
    crit: notable
    conf: 0.6
    for: [javascript]
    if:
      type: string
      regex: "setInterval.{0,50}axios"

  - id: setInterval-fetch
    desc: setInterval with fetch
    crit: notable
    conf: 0.6
    for: [javascript]
    if:
      type: string
      regex: "setInterval.{0,50}fetch"

  - id: setInterval-function
    desc: setInterval with function callback
    crit: notable
    conf: 0.7
    for: [javascript]
    if:
      type: string
      regex: "setInterval\\s*\\(\\s*function"

  - id: setTimeout-recursive
    desc: setTimeout in recursive pattern
    crit: notable
    conf: 0.7
    for: [javascript]
    if:
      type: string
      regex: "setTimeout.{0,50}setTimeout"

composite_rules:
  - id: setInterval-http
    desc: setInterval with HTTP calls (persistent
    crit: suspicious
    conf: 0.85
    attack: "T1053"
    for: [javascript]
    any:
      - id: setInterval-http-call
      - id: setInterval-axios
      - id: setInterval-fetch
