# objectives/persistence/startup/registry/windows.yaml
# Windows Startup Persistence Detection
# Detects Windows persistence mechanisms
# ATT&CK: T1547.001 (Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder)
defaults:
  attack: "T1547.001"
  crit: suspicious
  for: [pe, powershell, batch]
  platforms: [windows]
traits:
  # === Startup folder atomic indicators ===
  # Note: Rules use regex to match both raw strings (single backslash)
  # and source code strings (escaped double backslash)
  - id: start-menu-startup
    desc: Start Menu Programs Startup path
    crit: notable
    conf: 0.7
    if:
      type: string
      # Matches: \Start Menu\Programs\Startup or \\Start Menu\\Programs\\Startup
      regex: "\\\\{1,2}Start Menu\\\\{1,2}Programs\\\\{1,2}Startup"
  - id: appdata-startup
    desc: AppData Roaming Startup path
    crit: notable
    conf: 0.7
    if:
      type: string
      # Matches Roaming path leading to Startup folder
      regex: "Roaming.{1,50}Programs.{1,20}Startup"
  - id: shell-startup
    desc: shell:startup special folder
    crit: notable
    conf: 0.7
    if:
      type: string
      substr: "shell:startup"
      # === Registry Run key atomic indicators ===
  - id: software-run-key
    desc: SOFTWARE CurrentVersion Run path
    crit: notable
    conf: 0.7
    if:
      type: string
      substr: "SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run"
  - id: hkcu-run
    desc: HKCU Run key reference
    crit: notable
    conf: 0.6
    if:
      type: string
      regex: "HKCU.{0,50}Run"
  - id: hklm-run
    desc: HKLM Run key reference
    crit: notable
    conf: 0.6
    if:
      type: string
      regex: "HKLM.{0,50}Run"
  - id: currentversion-run
    desc: CurrentVersion Run path
    crit: notable
    conf: 0.7
    if:
      type: string
      substr: "CurrentVersion\\Run"
  - id: windows-registry-run-access
    desc: Accesses Windows registry Run key
    crit: notable
    conf: 0.75
    if:
      type: raw
      regex: "(?i)currentversion\\\\run"
  - id: currentversion-windows-key
    desc: CurrentVersion Windows key path
    crit: notable
    conf: 0.72
    if:
      type: string
      substr: "CurrentVersion\\Windows"
  - id: run-value-name
    desc: Run value name string
    crit: component
    conf: 0.6
    if:
      type: string
      exact: "Run"
      # === RunOnce atomic indicators ===
  - id: runonce-word
    desc: RunOnce keyword
    crit: notable
    conf: 0.5
    if:
      type: string
      regex: "(?i)(\\\\|/)RunOnce(\\\\|$)"
      # === Python winreg atomic indicators ===
  - id: import-winreg
    desc: import winreg statement
    crit: notable
    conf: 0.5
    if:
      type: string
      substr: "import winreg"
  - id: from-winreg
    desc: from winreg import
    crit: notable
    conf: 0.5
    if:
      type: string
      substr: "from winreg"
  - id: underscore-winreg
    desc: _winreg module (Python 2)
    crit: notable
    conf: 0.5
    if:
      type: string
      substr: "_winreg"
  - id: openkey-hkey
    desc: OpenKey HKEY call
    crit: notable
    conf: 0.6
    if:
      type: string
      regex: "OpenKey.{0,30}HKEY"
  - id: setvalueex
    desc: SetValueEx call
    crit: component
    conf: 0.6
    if:
      type: string
      substr: "SetValueEx"
  - id: createkey
    desc: CreateKey call
    crit: component
    conf: 0.6
    if:
      type: string
      substr: "CreateKey"
      # === Scheduled tasks atomic indicators ===
  - id: schtasks-create-cmd
    desc: schtasks /create command
    crit: notable
    conf: 0.7
    if:
      type: string
      regex: "schtasks\\s*/create"
  - id: schtasks-exe-create
    desc: schtasks.exe create
    crit: notable
    conf: 0.7
    if:
      type: string
      regex: "schtasks\\.exe.{0,30}create"
  - id: itaskscheduler-iface
    desc: ITaskScheduler interface
    crit: notable
    conf: 0.7
    if:
      type: string
      substr: "ITaskScheduler"
  - id: taskscheduler-word
    desc: TaskScheduler keyword
    crit: notable
    conf: 0.6
    if:
      type: string
      substr: "TaskScheduler"
      # === Service installation atomic indicators ===
  - id: sc-create-cmd
    desc: sc create command
    crit: notable
    conf: 0.7
    if:
      type: string
      regex: "sc\\s+create"
  - id: createservice-api
    desc: CreateService API
    crit: notable
    conf: 0.7
    if:
      type: string
      word: "CreateService"
  - id: installservice-fn
    desc: InstallService function
    crit: notable
    conf: 0.7
    if:
      type: string
      regex: "InstallService\\s*\\("
  - id: win32serviceutil-mod
    desc: win32serviceutil module
    crit: notable
    conf: 0.7
    if:
      type: string
      substr: "win32serviceutil"
      # === Shortcut file enumeration ===
  - id: lnk-endswith-filter
    desc: Filters files by .lnk extension
    crit: notable
    conf: 0.8
    if:
      type: string
      substr: ".lnk"
  - id: endswith-filter-content
    desc: Python endswith file extension filter
    crit: notable
    conf: 0.65
    if:
      type: raw
      substr: ".endswith("
      # === Windows shortcut manipulation ===
  - id: create-shortcut
    desc: COM CreateShortcut method call
    crit: notable
    conf: 0.75
    if:
      type: raw
      substr: "CreateShortcut"
  - id: shortcut-target-path
    desc: Shortcut TargetPath property access
    crit: notable
    conf: 0.75
    if:
      type: raw
      substr: ".TargetPath"
  - id: shortcut-arguments
    desc: Shortcut Arguments property modification
    crit: notable
    conf: 0.75
    if:
      type: raw
      substr: ".Arguments"
  - id: shortcut-save
    desc: Shortcut Save method call
    crit: notable
    conf: 0.65
    if:
      type: raw
      substr: ".Save()"
  - id: windows-powershell-registry
    desc: "Modifies registry for persistence via"
    conf: 0.9
    crit: suspicious
    if:
      type: trait
      id: micro-behaviors/os/registry/persistence::run-key-write-ps
composite_rules:
  # Startup folder composite
  - id: windows-startup-folder
    desc: "Targets Windows Startup folder for persistence"
    conf: 0.85
    crit: suspicious
    any:
      - id: start-menu-startup
      - id: appdata-startup
      - id: shell-startup
        # Registry Run key composite
  - id: windows-registry-run
    desc: "Accesses Windows registry Run key"
    conf: 0.9
    crit: suspicious
    any:
      - id: software-run-key
      - id: hkcu-run
      - id: hklm-run
      - id: currentversion-run
      - id: windows-registry-run-access
    downgrade:
      any:
        - id: metadata/quality::versioned
  - id: windows-registry-windows-run-value
    desc: Writes Run value in Windows key
    conf: 0.9
    crit: suspicious
    all:
      - id: currentversion-windows-key
      - id: run-value-name
      - id: micro-behaviors/os/registry/access::reg-set-value-ex-w
        # RunOnce composite
  - id: windows-registry-runonce
    desc: "Accesses Windows registry RunOnce key"
    conf: 0.9
    crit: suspicious
    any:
      - id: micro-behaviors/os/registry/persistence::run-once-key
      - id: runonce-word
        # Python winreg import composite
  - id: windows-python-winreg-import
    desc: "Imports Python winreg module for"
    conf: 0.7
    crit: notable
    any:
      - id: import-winreg
      - id: from-winreg
      - id: underscore-winreg
        # Python winreg composite
  - id: windows-python-winreg
    desc: "Accesses Windows registry via Python"
    conf: 0.8
    crit: notable
    # Stealth binaries are typically small
    size_max: 5000
    any:
      - id: import-winreg
      - id: from-winreg
      - id: underscore-winreg
      - id: openkey-hkey
      - id: setvalueex
      - id: createkey
    unless:
      - type: basename
        exact: "six.py"
        # Scheduled tasks composite
  - id: windows-schtasks
    desc: "Creates scheduled task for persistence"
    conf: 0.85
    crit: notable
    attack: "T1053.005"
    size_max: 500000
    any:
      - id: schtasks-create-cmd
      - id: schtasks-exe-create
      - id: itaskscheduler-iface
      - id: taskscheduler-word
        # Scheduled task create composite
  - id: windows-schtasks-create
    desc: "Creates Windows scheduled task"
    conf: 0.8
    crit: notable
    attack: "T1053.005"
    size_max: 500000
    any:
      - id: schtasks-create-cmd
      - id: schtasks-exe-create
      - id: itaskscheduler-iface
        # Service installation composite
  - id: windows-service-install
    desc: "Installs Windows service for persistence"
    conf: 0.85
    attack: "T1543.003"
    any:
      - id: sc-create-cmd
      - id: createservice-api
      - id: installservice-fn
      - id: win32serviceutil-mod
        # WMI persistence composite
  - id: windows-wmi-persistence
    desc: "Sets up WMI event subscription"
    conf: 0.9
    crit: suspicious
    attack: "T1546.003"
    any:
      - id: objectives/persistence/wmi/subscription::wmi-event-filter
      - id: objectives/persistence/wmi/subscription::wmi-commandline-consumer
      - id: objectives/persistence/wmi/subscription::wmi-activescript-consumer
        # Python registry persistence pattern
  - id: windows-python-reg-persistence
    desc: "Python Windows registry persistence"
    crit: suspicious
    conf: 0.95
    mbc: "F0006"
    all:
      - id: windows-python-winreg-import
      - id: windows-registry-run-access
        # Shortcut hijacking composite
  - id: windows-shortcut-hijack
    desc: "Windows shortcut hijacking via COM"
    crit: suspicious
    conf: 0.90
    attack: "T1547.009"
    all:
      - id: create-shortcut
      - id: shortcut-target-path
        # Shortcut argument injection composite
  - id: windows-shortcut-argument-injection
    desc: "Shortcut argument injection via COM"
    crit: suspicious
    conf: 0.90
    attack: "T1547.009"
    all:
      - id: create-shortcut
      - id: shortcut-arguments
      - id: shortcut-save
        # Shortcut enumeration + modification
        # Walks directories to find .lnk files and modifies them
  - id: windows-shortcut-walk-modify
    desc: "Enumerates and modifies shortcut files"
    crit: suspicious
    conf: 0.90
    attack: "T1547.009"
    all:
      - id: micro-behaviors/fs/directory/traverse::traverse-walk
      - id: create-shortcut
    any:
      - id: lnk-endswith-filter
      - id: endswith-filter-content
        # Native PE startup-folder resolution + shortcut filtering + registry writes
  - id: windows-native-shortcut-startup-tamper
    desc: Native startup shortcut and registry tamper chain
    crit: suspicious
    conf: 0.92
    all:
      - id: micro-behaviors/os/persistence/startup::startup-folder-shell-api-chain
      - id: lnk-endswith-filter
      - id: micro-behaviors/os/registry/access::registry-write-api-chain
    any:
      - id: micro-behaviors/os/registry/access::internet-settings-autodial-tamper
      - id: micro-behaviors/process/create/spawn::create-process-a
        # Startup folder persistence with registry and GUID mutex
  - id: pe-startup-folder-registry-persistence
    desc: Startup folder persistence with registry fingerprinting
    crit: suspicious
    conf: 0.9
    all:
      - id: micro-behaviors/os/persistence/startup::startup-folder-shell-path-composition-chain
      - id: micro-behaviors/os/registry/access::registry-write-api-chain
    any:
      - id: micro-behaviors/process/sync/mutex::guid-mutex
      - id: micro-behaviors/process/info/current::get-module-filename
      - id: micro-behaviors/process/create/spawn::create-process-w
    needs: 1
    downgrade:
      any:
        - id: metadata/signed/platform::microsoft
        - id: metadata/vendor/microsoft::is-microsoft-named
          # Full browser shortcut argument injection chain
  - id: windows-shortcut-walk-argument-inject
    desc: "Enumerates shortcuts and injects arguments"
    crit: hostile
    conf: 0.95
    attack: "T1547.009"
    mbc: "F0012"
    all:
      - id: micro-behaviors/fs/directory/traverse::traverse-walk
      - id: create-shortcut
      - id: shortcut-arguments
      - id: shortcut-save
        # Scheduled task with suspicious payload
  - id: windows-schtasks-payload
    desc: "Scheduled task with executable payload"
    crit: suspicious
    conf: 0.95
    attack: "T1053.005"
    mbc: "F0006"
    all:
      - id: windows-schtasks-create
      - id: micro-behaviors/process/create/shell::cmd-exe
