# objectives/persistence/login/registry/autostart/traits.yaml
defaults:
  for: [python, javascript, vbs, powershell, pe, c, cpp]
  crit: notable
  conf: 0.7
  platforms:
    - windows
traits:
  - id: run-key
    desc: Registry Run key
    if:
      type: text
      regex: '(?i)Software[\\/]Microsoft[\\/]Windows[\\/]CurrentVersion[\\/]Run([\\/]|$|\s)'
    unless:
      - id: well-known/tool/detection/malware-config-parser
      - id: metadata/binary/installer/family/
      - id: objectives/command-and-control/dropper/execution/resource/
      - type: text
        exact: WMI Win32 Provider
        section: .rsrc
      - type: text
        regex: (?i)[\\/]var[\\/]+run\b
  - id: currentversion-run-key-path
    desc: CurrentVersion Run registry path
    crit: notable
    conf: 0.9
    if:
      type: text
      substr: 'Software\Microsoft\Windows\CurrentVersion\Run'
    unless:
      - id: metadata/binary/installer/family/
      - id: metadata/signed/platform/
      - id: well-known/app/peazip/
  - id: run-once-key
    desc: Registry RunOnce key
    if:
      type: text
      regex: (?i)\\RunOnce\b
    unless:
      - id: objectives/evasion/self-delete/cleanup/heartbeat::advanced-inf-setup-regpath
      - id: metadata/signed/platform::microsoft-corporation
      - id: metadata/binary/installer/family/
  - id: run-service-key
    desc: Registry RunServices key
    if:
      type: text
      regex: (?i)\\RunServices\b
    unless:
      - type: text
        exact: WMI Win32 Provider
        section: .rsrc
  - id: winlogon-shell-key
    desc: Registry Winlogon Shell key
    if:
      type: text
      regex: (?i)\\Winlogon\\Shell\b
  - id: userinit-key
    desc: Registry Userinit key
    if:
      type: text
      regex: (?i)\\Winlogon\\Userinit\b
  - id: set-itemproperty-run
    desc: Set-ItemProperty Run
    if:
      type: text
      regex: Set-ItemProperty.{0,50}Run
  - id: new-itemproperty-run
    desc: New-ItemProperty Run
    if:
      type: text
      regex: '(?i)New-ItemProperty.{0,120}CurrentVersion\\Run\b'
  - id: powershell-hkcu-run-path
    desc: PowerShell HKCU Run key path
    crit: notable
    conf: 0.96
    for: [powershell]
    if:
      type: text
      substr: 'HKCU:\Software\Microsoft\Windows\CurrentVersion\Run'
  - id: powershell-newitemproperty-regpath
    desc: New-ItemProperty writes RegPath variable
    crit: component
    conf: 0.95
    for: [powershell]
    if:
      type: text
      regex: '(?i)New-ItemProperty\s+-Path\s+\$RegPath'
  - id: powershell-appdata-exe-target
    desc: PowerShell builds an AppData executable path
    crit: component
    conf: 0.94
    for: [powershell]
    if:
      type: text
      regex: '(?i)\$TargetPath\s*=\s*Join-Path\s+\$env:APPDATA.{1,160}\.exe'
  - id: reg-add-run
    desc: reg add Run
    if:
      type: text
      regex: (?i)reg\s+add.{0,80}\\Run\b
  - id: quoted-programdata-prefix
    desc: Quoted ProgramData path prefix
    crit: component
    conf: 0.9
    for: [pe, dll]
    if:
      type: text
      substr: '"C:\ProgramData\'
  - id: quoted-exe-double-dash-arg
    desc: Quoted executable with double-dash arg
    crit: component
    conf: 0.88
    for: [pe, dll]
    if:
      type: text
      regex: '(?i)\.exe"\s+--[A-Z0-9]{2,12}'
  - id: short-camel-on-run-value-token
    desc: Short CamelCase token adjacent to Run key
    crit: component
    conf: 0.78
    for: [pe, dll]
    if:
      type: text
      regex: '(?i)(\\Run|CurrentVersion).{0,96}\b[A-Z][A-Za-z0-9]{2,8}On\b'
    unless:
    - type: trait
      id: metadata/file/extension/identity::linux-kernel-package-artifact-basename
  - id: encoded-run-key-path
    desc: Encoded Run registry path
    crit: component
    conf: 0.98
    for: [pe]
    if:
      type: text
      exact: 0Nzl09Te0eLL2vYA8fzw/AXzy9T2/QP89PDL4PLx8QL989UC8fD2/P3L0fL9nw==
  - id: encoded-geekbuddy-run-name
    desc: Encoded GeekBuddy Run value name
    crit: component
    conf: 0.98
    for: [pe]
    if:
      type: text
      exact: 0PP+8fO/5AIC+OHyAwPmnw==
  - id: batch-reg-add-run
    desc: Batch reg add Run
    crit: component
    conf: 0.9
    for: [batch]
    if:
      type: text
      substr: '%rgs%g add'
  - id: batch-run-key-fragment
    desc: Batch Run key fragment
    crit: component
    conf: 0.9
    for: [batch]
    if:
      type: text
      substr: '\R%u%n'
  - id: batch-public-windowssec-path
    desc: Batch Public WindowsSec path
    crit: component
    conf: 0.92
    for: [batch]
    if:
      type: text
      substr: WindowsSec
  - id: batch-ntoskrnl-masquerade
    desc: Batch ntoskrnl masquerade
    crit: component
    conf: 0.92
    for: [batch]
    if:
      type: text
      substr: 'ntoskrnl%eA%'
composite_rules:
  - id: powershell-run-path-variable-write
    desc: PowerShell writes a variable-backed Run key
    crit: suspicious
    conf: 0.96
    for: [powershell]
    near_lines: 20
    all:
      - id: powershell-hkcu-run-path
      - id: powershell-newitemproperty-regpath
  - id: common-persistence-keys
    desc: Common Windows registry persistence keys
    any:
      - id: run-key
      - id: currentversion-run-key-path
      - id: run-once-key
      - id: run-service-key
      - id: winlogon-shell-key
      - id: userinit-key
  - id: run-key-write-ps
    desc: Modifies Run key via PowerShell
    any:
      - id: set-itemproperty-run
      - id: new-itemproperty-run
      - id: reg-add-run
  - id: quoted-programdata-exe-run-command
    desc: Quoted ProgramData executable Run command
    crit: component
    conf: 0.9
    for: [pe, dll]
    near_bytes: 256
    all:
      - id: quoted-programdata-prefix
      - id: quoted-exe-double-dash-arg
  - id: run-key-near-candidate-value-name
    desc: Run key with nearby candidate value name
    crit: notable
    conf: 0.82
    for: [pe, dll]
    near_bytes: 384
    all:
      - id: run-key
      - id: short-camel-on-run-value-token
    any:
      - id: quoted-exe-double-dash-arg
      - id: micro-behaviors/fs/path/extension::exe-dup
  - id: programdata-run-key-command
    desc: ProgramData executable persisted through Run key
    crit: suspicious
    conf: 0.92
    for: [pe, dll]
    all:
      - id: run-key
      - id: quoted-programdata-exe-run-command
      - id: metadata/signed/trust-level::unsigned
    unless:
      - id: metadata/binary/installer/family/
      - id: metadata/signed/platform/
  - id: unsigned-microsoft-run-key-persistence
    desc: Unsigned Microsoft-claim PE carries Run-key persistence
    crit: hostile
    conf: 0.96
    for: [pe, dll]
    attack: T1547.001
    all:
      - id: objectives/evasion/masquerade/trust::unsigned-pe-microsoft-version-claim
      - id: currentversion-run-key-path
    any:
      - id: objectives/anti-static/obfuscation/payload/registry::reg-open-key
      - id: objectives/persistence/login/registry/run-key::currentversion-run-key-literal
    unless:
      - id: metadata/binary/installer/family/
      - id: metadata/signed/platform/
      - id: well-known/app/microsoft-office-c2r::skype-business-recording-manager
  - id: encoded-run-key-persistence
    desc: Encoded Run key persistence
    crit: suspicious
    conf: 0.98
    for: [pe]
    all:
      - id: encoded-run-key-path
      - id: encoded-geekbuddy-run-name
      - id: micro-behaviors/os/registry/manipulate::registry-write-api-chain
    any:
      - id: objectives/anti-static/obfuscation/string/runtime-decrypt::base64-arithmetic-xor-decoder
      - id: micro-behaviors/os/registry/access::reg-set-value-ex-a
  - id: batch-public-windowssec-run-key
    desc: Batch Public WindowsSec Run key
    crit: suspicious
    conf: 0.96
    for: [batch]
    all:
      - id: batch-reg-add-run
      - id: batch-run-key-fragment
      - id: batch-public-windowssec-path
      - id: batch-ntoskrnl-masquerade
  - id: powershell-appdata-run-key-payload
    desc: PowerShell drops, persists, and starts an AppData executable
    crit: hostile
    conf: 0.99
    for: [powershell]
    attack: T1547.001
    all:
      - id: powershell-run-path-variable-write
      - id: powershell-appdata-exe-target
      - id: micro-behaviors/fs/write/file/direct::write-all-bytes
      - id: micro-behaviors/process/create/start-process::start-process-call
  - id: scrambled-upx-run-key-dropper
    desc: Tampered UPX binary persists through a Run key and spawns a process
    crit: hostile
    conf: 0.99
    for: [pe]
    attack: T1027.002,T1547.001,T1106
    all:
      - id: objectives/anti-static/pack/upx::scrambled-upx-header
      - id: metadata/binary/layout::high-entropy-overlay
      - id: run-key
      - id: micro-behaviors/process/create/spawn::create-process-w
