# objectives/lateral-movement/supply-chain/vscode/malicious.yaml
# Malicious VSCode Extension Patterns
# Detects high-confidence malicious patterns in VSCode extensions
# These combine multiple indicators to reduce false positives
defaults:
  platforms: [linux, macos, windows, unix, android, ios]
  for: [scripts, source, documents]
traits:
  # NOTE: YARA kept for api-usage - requires fullword match on "activate" which
  # cleave doesn't support, plus regex patterns for require/import
  - id: api-usage
    desc: Uses VSCode API (extension marker)
    crit: notable
    conf: 1.0
    for: [javascript, typescript]
    if:
      id: micro-behaviors/process/create/shell/host::vscode-extension-marker
  - id: postInstall-word
    desc: postInstall keyword
    crit: notable
    conf: 0.7
    for: [javascript, typescript]
    if:
      type: string_value
      substr: "postInstall"
  - id: exec-pattern
    desc: exec function pattern
    crit: notable
    conf: 0.7
    for: [javascript, typescript]
    if:
      type: raw
      regex: "exec(Sync|Async)\\("
      # === Atomic indicators for credential-theft ===
      # Removed ssh-dir duplicate - use micro-behaviors/fs/path/sensitive/ssh::keys
  - id: readFile-pattern
    # aws-dir: use objectives/credential-access/cloud/token::aws-config-dup (canonical)

    # Removed npmrc-file duplicate - use objectives/credential-access/environment/secret-access::dot-npmrc
    desc: readFile function pattern
    crit: component
    conf: 0.6
    for: [javascript, typescript]
    if:
      type: string_value
      regex: "read(File|FileSync)"
      # VSCode namespace usage pattern
  - id: vscode-namespace
    desc: Uses VSCode namespace
    crit: notable
    conf: 1.0
    for: [javascript, typescript]
    if:
      id: micro-behaviors/process/create/shell/host::namespace
      # File copy operations
  - id: file-copy-ops
    desc: File copy operations
    crit: notable
    conf: 1.0
    for: [javascript, typescript]
    if:
      type: string_value
      regex: "copy(File|FileSync)"
composite_rules:
  # State helpers
  - id: state-markers
    desc: State management markers (globalState or postInstall)
    crit: notable
    conf: 0.75
    for: [javascript, typescript]
    any:
      - id: micro-behaviors/process/create/shell/host::global-state
      - id: postInstall-word
        # Credential paths
  - id: credential-paths
    desc: Sensitive credential file/directory paths
    crit: notable
    conf: 0.8
    for: [javascript, typescript]
    any:
      - id: micro-behaviors/fs/path/sensitive/ssh::keys
      - id: objectives/credential-access/cloud/token::aws-config-dup
      - id: objectives/credential-access/environment/secret-access::dot-npmrc
      - id: metadata/quality/config::dotenv-file
        # Send methods
  - id: send-methods
    desc: HTTP send methods
    crit: notable
    conf: 0.7
    for: [javascript, typescript]
    any:
      - id: micro-behaviors/communications/http/post/
      - id: micro-behaviors/communications/http/request/
      - id: micro-behaviors/communications/http/request::axios
        # Post-Install Code Execution
  - id: post-install-exec
    desc: Post-install script execution (dropper pattern)
    crit: suspicious
    conf: 0.95
    mbc: "B0024"
    attack: "T1059"
    for: [javascript, typescript]
    all:
      - id: objectives/lateral-movement/supply-chain/vscode::function-activate
      - id: state-markers
      - id: exec-pattern
        # Extension Credential Theft
  - id: credential-theft-comp
    desc: VSCode extension accessing sensitive credentials
    crit: suspicious
    conf: 0.9
    mbc: "E1519"
    attack: "T1555"
    for: [javascript, typescript]
    all:
      - id: credential-paths
      - id: readFile-pattern
      - id: send-methods
        # Staged Decryption and Execution inside a VSCode Extension
  - id: persistence-locations
    #  - id: staged-dropper-comp
    #    desc: VSCode extension with staged decryption
    #    crit: suspicious
    #    conf: 0.99
    #    mbc: "B0032"
    #    attack: "T1027"
    #    for: [javascript, typescript]
    #    all:
    #      - id: objectives/command-and-control/dropper/encrypted::js-decrypt-eval
    #      - id: api-usage
    #      - id: vscode-namespace

    desc: References persistence locations (startup/autorun paths)
    crit: notable
    conf: 0.8
    for: [javascript, typescript]
    any:
      - id: micro-behaviors/fs/path/persistence-locations/keywords::startup-keyword
      - id: micro-behaviors/os/schedule/cron/keywords::crontab-keyword
      - id: micro-behaviors/fs/path/persistence-locations/keywords::launchd-keyword
      - id: objectives/persistence/systemd/service::systemd-keyword
        # Extension Persistence via well-known hooks
  - id: persistence
    desc: Extension establishing persistence
    crit: suspicious
    conf: 0.95
    mbc: "F0011"
    attack: "T1547"
    for: [javascript, typescript]
    all:
      - id: api-usage
      - id: persistence-locations
      - id: micro-behaviors/fs/write/file::js-file-write
  - id: file-write-ops
    desc: File write or copy operations
    crit: notable
    conf: 1.0
    for: [javascript, typescript]
    any:
      - id: micro-behaviors/fs/write/file::js-file-write
      - id: file-copy-ops
