# objectives/anti-static/obfuscation/encoding/base64/traits.yaml
# Encoded String Obfuscation Detection
# Detects encoded sensitive paths and credentials (base64, xor, hex, etc.)
# ATT&CK: T1027.001 (Obfuscation: Steganography or Compression)
# MBC: F0001 (Obfuscation)

traits:
  # === Encoded Persistence Paths (any encoding: base64, xor, hex, etc.) ===
  - id: base64-systemd-path
    desc: "Encoded systemd path (base64/xor/hex)"
    crit: suspicious
    conf: 0.95
    for: [elf, shell]
    if:
      type: encoded
      substr: "/usr/lib/systemd/systemd"
    attack: "T1547"

  - id: base64-ssh-sftp-path
    desc: "Encoded SSH SFTP path (base64/xor/hex)"
    crit: suspicious
    conf: 0.95
    for: [elf, shell]
    if:
      type: encoded
      substr: "/usr/libexecution/openssh/sftp-server"
    attack: "T1021.004"

  # === Encoded DVR Targeting Paths (any encoding: base64, xor, hex, etc.) ===
  - id: base64-hikvision-firmware-path
    desc: "Encoded Hikvision path (base64/xor/hex)"
    crit: suspicious
    conf: 0.95
    for: [elf, shell]
    if:
      type: encoded
      substr: "usr/dvr_main _8182T_1108"
    attack: "T1046"

  - id: base64-anko-dvr-path
    desc: "Encoded Anko DVR path (base64/xor/hex)"
    crit: suspicious
    conf: 0.95
    for: [elf, shell]
    if:
      type: encoded
      substr: "anko-app/ankosample _8182T_1104"
    attack: "T1046"

  # === Generic Base64 Encoding Pattern ===
  - id: base64-path-prefix
    desc: "Base64-encoded system path"
    crit: suspicious
    conf: 0.8
    for: [elf, shell]
    size_max: 50000
    count_min: 3
    if:
      type: encoded
      encoding: base64
      regex: "^/[a-z]{2,}[a-z_/]*$"
    attack: "T1027.001"

  # Note: Base64-like strings are common in legitimate binaries (certificates, keys)
  # This pattern alone is not suspicious - requires context or high density
  - id: base64-suspicious-encoding
    desc: "Suspicious base64 patterns"
    crit: notable
    conf: 0.5
    for: [elf]
    if:
      type: string
      regex: "^[a-zA-Z0-9+/]{20,}={1,3}$"
    not:
      # Exclude standard base64 alphabet lookup table
      - "ABCDEFGHIJKLMNOPQRSTUVWXYZabcdefghijklmnopqrstuvwxyz0123456789+/="
    attack: "T1027.001"

  # Dense base64 strings indicate heavy encoding of payloads or configs.
  # Legitimate binaries rarely have more than a few base64 strings; malware
  # with encoded payloads or C2 configs often has 5+ dense base64 blobs.
  - id: dense-base64-encoding
    desc: "Dense base64-encoded content"
    crit: suspicious
    conf: 0.85
    for: [elf, shell]
    size_max: 200000
    count_min: 5
    per_kb_min: 0.3
    if:
      type: string
      regex: "^[a-zA-Z0-9+/]{40,}={0,2}$"
    not:
      - "ABCDEFGHIJKLMNOPQRSTUVWXYZabcdefghijklmnopqrstuvwxyz0123456789+/="
      - "RgA1ADAANwA4AEYAMwA1AC0AQwA1ADUAMQAtADEAMQBEADMALQA4ADkAQgA5AC0AMAAwADAAMABGADgAMQBGAEUAMgAyADEA"
      - "nuwmopywnwcjmihwauaevgkbimtnfaiuweikwggjdlupvcnvypcuymxgzsvgudarkidgrkrwvmnqqcoqsigP"
      - "6879706572766563746f722e6d652e6476656465762e7275"
      - "bb1131f893f91f1bf5461285b26c0b622d21a37e"
    attack: "T1027.001"

composite_rules:
  - id: base64-encoded-dvr-targeting
    desc: "Encoded DVR targeting (base64/xor/hex)"
    crit: suspicious
    conf: 0.93
    for: [elf, shell]
    attack: "T1046"
    any:
      - id: base64-hikvision-firmware-path
      - id: base64-anko-dvr-path

  - id: base64-obfuscation-evidence
    desc: "Base64 obfuscation evidence"
    crit: suspicious
    conf: 0.85
    for: [elf, shell]
    attack: "T1027.001"
    mbc: "F0001"
    all:
      - id: base64-path-prefix
      - id: base64-suspicious-encoding

  - id: mirai-base64-obfuscation
    desc: "Mirai base64 obfuscation"
    crit: suspicious
    conf: 0.92
    for: [elf, shell]
    attack: "T1027.001"
    mbc: "F0001"
    all:
      - id: well-known/malware/botnet/mirai
      - id: base64-obfuscation-evidence

