# objectives/persistence/system/cron/job/linux.yaml
defaults:
  attack: T1053.003
  mbc: B0028
  platforms:
    - unix
    - windows
  for: [scripts, binaries]
traits:
  - id: etc-cron-path
    desc: /etc/cron system cron path
    crit: suspicious
    conf: 0.7
    if:
      type: text
      regex: \b/etc/cron(/|$)
    unless:
      - id: metadata/file/format/
      - id: metadata/file/catalog/identity/
  - id: cron-d-path-reference
    desc: /etc/cron.d system cron path
    crit: notable
    conf: 0.7
    if:
      type: text
      substr: /etc/cron.d
    unless:
      - id: well-known/tool/detection/security-scanner/
      - id: well-known/tool/detection/security-scanner::wpsecscan-cli-context
      - id: metadata/package/quality/false-positives::benign-app-fp-context
      - id: metadata/build/package::pacman-post-remove-cron-cleanup
      - id: metadata/build/package::debian-postrm-cron-cleanup
      - id: well-known/app/olares::cli-source-tree
      - type: text
        substr: CloudWatchMetrics
      - type: text
        regex: /etc/cron\.d/(glusterfs_snap_cron_tasks|gcron_update_task)
      - type: text
        substr: PERSISTENCE_PATHS
  - id: crontab-at-reboot
    desc: Crontab @reboot directive
    crit: suspicious
    conf: 0.9
    for: [go, elf, macho, shell, python]
    if:
      # substr, not word: '@' is non-word so word-boundary matchers miss @reboot
      type: text
      substr: '@reboot'
    unless:
      - id: metadata/package/quality/recognized-contexts::malcontent-defensive-scanner-context
  - id: crontab-at-reboot-root-payload
    desc: Cron @reboot root payload template
    crit: suspicious
    conf: 0.95
    for: [elf, shell]
    platforms: [linux, android]
    if:
      type: encoded
      encoding: xor
      substr: '@reboot root'
    unless:
      - id: metadata/package/quality/recognized-contexts::malcontent-defensive-scanner-context
  - id: crontab-edit
    desc: Crontab edit operation
    crit: notable
    conf: 0.8
    for: [go, elf, macho, shell, python]
    if:
      type: text
      regex: crontab.{0,10}-e
    unless:
      - id: well-known/tool/detection/security-scanner/
      - id: well-known/app/olares::cli-source-tree
composite_rules:
  - id: cron
    desc: Cron job persistence indicator
    crit: suspicious
    conf: 0.66
    any:
      - id: micro-behaviors/os/autorun/scheduled::crontab-keyword
      - id: etc-cron-path
      - id: micro-behaviors/fs/path/config/system::var-spool-cron-path
      - id: objectives/persistence/system/cron/job::cron-d-path-reference
    unless:
      - id: well-known/app/kilocode::kilo-code-native-runtime
      - id: well-known/app/reg-devshell::developer-shell-application
  - id: crontab-reboot
    desc: Crontab reboot persistence
    crit: suspicious
    conf: 0.9
    for: [go, elf, macho, shell, python]
    any:
      - id: crontab-at-reboot
      - id: crontab-at-reboot-root-payload
      - id: crontab-edit
    unless:
      - id: metadata/package/quality/recognized-contexts::malcontent-defensive-scanner-context
  - id: cron-source-query-ddos-bot
    desc: Cron reboot persistence with Source-query DDoS payload
    crit: hostile
    conf: 0.98
    attack: T1053.003,T1498.001
    mbc: B0030
    for: [elf]
    platforms: [linux, android]
    all:
      - id: crontab-reboot
      - id: objectives/impact/dos/attack/botnet::source-engine-query-payload
    any:
      - id: crontab-at-reboot-root-payload
      - id: objectives/persistence/system/service/systemd::native-unit-authoring
    unless:
      - id: metadata/package/quality/recognized-contexts::malcontent-defensive-scanner-context
  - id: cron-preload-rootkit-persistence
    desc: Cron persists LD_PRELOAD rootkit
    crit: hostile
    conf: 0.98
    for: [elf, shell]
    all:
      - id: cron
      - id: objectives/evasion/hijack-execution-flow/preload::ld-preload-rootkit-persistence
