# objectives/persistence/system/cron/job/python.yaml
defaults:
  platforms: [unix]
  for: [python, pyc]
  attack: T1053.003
traits:
  - id: python-self-cron-entry
    desc: Python cron entry runs current file
    crit: suspicious
    conf: 0.9
    for: [python]
    if:
      type: text
      regex: '(__file__|script_path).{0,120}cron_entry|cron_entry.{0,120}(__file__|script_path)'
  - id: python-pyc-self-cron-entry
    desc: Python bytecode self cron entry
    crit: suspicious
    conf: 0.88
    for: [pyc]
    # Anchor on the cron_entry token (same discriminator as the source-level
    # python-self-cron-entry); bare `__file__` is in nearly every module.
    if:
      type: text
      substr: cron_entry
  - id: python-health-check-cron-script
    desc: Cron entry references health-check script
    crit: suspicious
    conf: 0.9
    if:
      type: text
      regex: 'crontab.{0,120}health_check\.sh|health_check\.sh.{0,120}crontab'
    unless:
      - id: metadata/package/testing/
  - id: python-crontab-beacon-callback
    desc: Crontab install with beacon callback
    crit: suspicious
    conf: 0.92
    if:
      type: text
      regex: 'crontab\s+-[\s\S]{0,220}(heartbeat|keepalive|beacon)'
    unless:
      - id: well-known/tool/detection/llm-scanner::mvar-security-context
composite_rules:
  - id: python-self-cron-persistence
    desc: Python self crontab persistence
    crit: hostile
    conf: 0.96
    all:
      - id: micro-behaviors/os/autorun/scheduled::crontab-install-command
    any:
      - id: python-self-cron-entry
      - id: python-pyc-self-cron-entry
    needs: 1
  - id: python-shell-pipe-cron-persistence
    desc: Python shell-piped cron persistence
    crit: hostile
    conf: 0.96
    all:
      - id: micro-behaviors/os/autorun/scheduled::crontab-shell-pipe
      - id: python-health-check-cron-script
  - id: python-pyc-health-cron-persistence
    desc: Python bytecode health cron persistence
    crit: hostile
    conf: 0.94
    for: [pyc]
    all:
      - id: python-health-check-cron-script
      - id: python-pyc-self-cron-entry
  - id: python-telegram-agent-cron-persistence
    desc: Telegram remote agent with crontab persistence
    crit: hostile
    conf: 0.99
    scope: archive
    for: [whl]
    attack: T1053.003,T1059.006,T1102
    all:
      - id: micro-behaviors/communications/http/services/telegram::telegram-api-host
      - id: micro-behaviors/communications/http/services/telegram::python-telegram-get-updates
      - id: objectives/command-and-control/backdoor/control/telegram::python-telegram-input-queue
      - id: micro-behaviors/process/create/subprocess::variable-command-shell-true
      - id: micro-behaviors/os/autorun/scheduled::crontab-install-command
