# objectives/persistence/system/init/boot/linux.yaml
defaults:
  for:
    - elf
    - shell
    - python
    - javascript
    - typescript
    - go
    - c
    - cpp
    - php
  attack: T1547.013
  crit: suspicious
  platforms:
    - linux
traits:
  - id: config-autostart
    desc: .config/autostart directory
    crit: notable
    conf: 0.7
    if:
      type: text
      regex: '\.config.{0,40}autostart'
    unless:
      - id: metadata/package/quality/false-positives-scanner::installtime-audit-defensive-scanner-context
  - id: etc-xdg-autostart
    desc: /etc/xdg/autostart directory
    crit: suspicious
    conf: 0.7
    if:
      type: text
      substr: /etc/xdg/autostart
    unless:
      - id: metadata/file/profile/identity::generated-attack-map-context
      - id: well-known/app/openbox::package-code
      - id: well-known/app/open-vm-tools::debian-packaging
      - id: well-known/app/clonezilla::package-context
      - id: well-known/app/clonezilla::archive-name
      - id: well-known/app/clonezilla::clonezilla-script-context
  - id: xdg-config-autostart
    desc: XDG autostart path literal
    crit: suspicious
    conf: 0.6
    unless:
      - id: well-known/app/telegram::telegram-desktop-third-party-source
      - id: well-known/app/openbox::package-code
      - id: well-known/app/open-vm-tools::debian-packaging
    if:
      type: literal
      regex: XDG_CONFIG_HOME.{0,50}autostart
  - id: desktop-entry-header
    desc: '[Desktop Entry] header'
    crit: component
    conf: 0.6
    if:
      type: text
      substr: '[Desktop Entry]'
  - id: x-gnome-autostart
    desc: X-GNOME-Autostart directive
    crit: notable
    conf: 0.7
    if:
      type: text
      substr: X-GNOME-Autostart
  - id: local-share-systemd-user
    desc: .local/share/systemd/user directory
    crit: notable
    conf: 0.7
    if:
      type: text
      substr: .local/share/systemd/user
  - id: systemctl-user-enable
    desc: systemctl --user enable
    crit: notable
    conf: 0.7
    if:
      type: text
      regex: '\bsystemctl\b[^\r\n]{0,30}\s--user\s+enable\b'
    unless:
      - id: well-known/tool/detection/security-scanner/
      # AI-agent observability tools install their own telemetry collector daemon.
      - id: well-known/tool/detection/ai-guardrail/
      - id: well-known/app/calvyn-code
      - id: metadata/package/quality/recognized-contexts::pvetui-workflow-test-script
  - id: etc-rc-local
    desc: /etc/rc.local file reference
    crit: notable
    conf: 0.8
    attack: T1037.004
    if:
      type: text
      substr: /etc/rc.local
    unless:
      - type: text
        substr: PERSISTENCE_PATHS
  - id: etc-rc-d-rc-local
    desc: /etc/rc.d/rc.local file
    crit: notable
    conf: 0.8
    attack: T1037.004
    if:
      type: text
      substr: /etc/rc.d/rc.local
    unless:
      - type: text
        substr: PERSISTENCE_PATHS
  - id: openrc-init-script-header
    desc: OpenRC init script header
    crit: notable
    conf: 0.85
    attack: T1037.004
    for: [shell]
    if:
      type: text
      substr: "#!/sbin/openrc-run"
  - id: update-rc-d-cmd
    desc: update-rc.d command
    crit: notable
    conf: 0.7
    attack: T1037.004
    if:
      type: text
      substr: update-rc.d
    unless:
      - id: well-known/app/dropbear::source-tree
  - id: chkconfig-cmd
    desc: chkconfig command (persistence context)
    crit: suspicious
    conf: 0.7
    attack: T1037.004
    if:
      type: text
      substr: chkconfig --add
    unless:
      - id: metadata/build/package/
      - id: well-known/tool/dual-use/shadowsocks-libev::shadowsocks-libev-source
      - id: metadata/package/files/system-package::os-package-manager-tool
      - id: well-known/app/walinuxagent::source-tree
  - id: initd-path-before-service-action
    desc: init.d path before service action
    crit: notable
    conf: 0.85
    attack: T1037.004
    if:
      type: text
      regex: /etc/init\.d/[^\s"'`|;]{1,96}.{0,20}\s(start|stop|restart)\b
    unless:
      - id: well-known/app/dropbear::source-tree
    not:
      - substr: complete -x -p "/etc/init.d/*"
      - substr: -a start --description
      - substr: -a restart --description
      - regex: /etc/init\.d/\*" -a (start|restart) --description
  - id: initd-path-before-install-action
    desc: init.d path before install action
    crit: suspicious
    conf: 0.85
    attack: T1037.004
    if:
      type: text
      regex: /etc/init\.d/[^\s"'`|;]{1,96}.{0,40}\b(install|enable|disable)\b
    unless:
      - id: well-known/app/walinuxagent::source-tree
  - id: install-action-before-initd-path
    desc: Install action before init.d path
    crit: notable
    conf: 0.85
    attack: T1037.004
    unless:
      - id: objectives/persistence/system/init/boot/
      - id: well-known/app/dropbear::source-tree
    if:
      type: text
      regex: \b(install|enable|disable)\b.{0,40}/etc/init\.d/[^\s"'`|;]{1,96}
  - id: write-action-before-initd-path
    desc: Write action before init.d path
    crit: suspicious
    conf: 0.85
    attack: T1037.004
    unless:
      - id: metadata/build/docker::dockerfile-basename
      - id: objectives/persistence/system/init/boot/
      - id: metadata/build/package/
    if:
      type: text
      regex: \b(copy|write|chmod|chown)\b.{0,40}/etc/init\.d/[^\s"'`|;]{1,96}
  - id: exec-eq-directive
    desc: Exec= directive
    crit: component
    conf: 0.6
    if:
      type: text
      substr: Exec=
    unless:
    - type: trait
      id: metadata/file/extension/identity::linux-kernel-package-artifact-basename
  - id: hidden-false-directive
    desc: Hidden=false directive
    crit: component
    conf: 0.5
    if:
      type: text
      substr: Hidden=false
  - id: type-application-directive
    desc: Type=Application directive
    crit: component
    conf: 0.5
    if:
      type: text
      substr: Type=Application
  - id: mattermost-desktop-resource-marker
    desc: Mattermost Desktop resource marker
    crit: notable
    conf: 0.95
    if:
      type: text
      regex: '(Mattermost Desktop|mattermost-desktop|app\.mattermost\.desktop)'
composite_rules:
  - id: linux-xdg-autostart
    desc: XDG autostart directory access
    conf: 0.85
    crit: suspicious
    any:
      - id: config-autostart
      - id: etc-xdg-autostart
      - id: xdg-config-autostart
    unless:
      - id: metadata/package/quality/recognized-contexts::malcontent-defensive-scanner-context
      - id: well-known/app/clonezilla::package-context
      - id: well-known/app/clonezilla::clonezilla-script-context
  - id: linux-desktop-entry
    desc: Desktop entry file creation (.desktop)
    conf: 0.8
    crit: notable
    any:
      - id: desktop-entry-header
      - id: x-gnome-autostart
  - id: linux-user-systemd
    desc: User-level systemd service persistence
    conf: 0.85
    crit: notable
    attack: T1543.002
    for: [elf, shell, python]
    any:
      - id: objectives/persistence/system/service/systemd::systemd-user-service-path
      - id: local-share-systemd-user
      - id: systemctl-user-enable
  - id: linux-rc-local
    desc: rc.local persistence mechanism
    conf: 0.9
    crit: suspicious
    attack: T1037.004
    any:
      - id: etc-rc-local
      - id: etc-rc-d-rc-local
  - id: linux-init-script
    desc: SysV init script with installation or enablement action
    conf: 0.85
    crit: suspicious
    attack: T1037.004
    unless:
      - id: well-known/app/openssh/
      - id: well-known/app/frr::frr-init-script-context
      - id: well-known/app/tailscale::tailscale-init-script-context
    all:
      - id: micro-behaviors/os/service/config::sysv-init-script-header
    any:
      - id: update-rc-d-cmd
      - id: chkconfig-cmd
      - id: initd-path-before-service-action
      - id: initd-path-before-install-action
      - id: install-action-before-initd-path
      - id: write-action-before-initd-path
  - id: linux-init-iot-bot-persistence
    desc: Linux init persists IoT bot
    conf: 0.98
    crit: hostile
    attack: T1037.004
    for: [elf, shell]
    all:
      - id: linux-init-script
      - id: objectives/command-and-control/botnet/iot::huawei-upnp-encoded-download-exec

  - id: multi-method-elf-persistence-templates
    desc: ELF embeds cron, systemd, and init persistence templates
    crit: hostile
    conf: 0.98
    attack: T1053.003,T1543.002,T1037.004
    for: [elf]
    all:
      - id: objectives/persistence/system/init/boot::native-systemd-unit-path
      - id: micro-behaviors/os/service/config::init-script-path
      - id: micro-behaviors/os/service/config::sysv-init-script-header
      - id: objectives/persistence/system/cron/schedule::crontab-echo-pipe
  - id: debhelper-installinit-maintainer-block
    desc: Debhelper init maintainer block
    crit: notable
    conf: 0.95
    all:
      - id: metadata/build/package::debhelper-dh-installinit-marker
      - id: metadata/build/package::debhelper-auto-section-end
  - id: linux-xdg-autostart-dir
    desc: XDG autostart directory path
    conf: 0.8
    crit: notable
    any:
      - id: config-autostart
      - id: etc-xdg-autostart
    unless:
      - id: well-known/app/macos-benign::mattermost-desktop-bundle
      - id: well-known/app/macos-benign::mattermost-asar-resource
      - id: mattermost-desktop-resource-marker
      - id: well-known/app/bitwarden::bitwarden-text-marker
      - id: well-known/app/rclone-wiz::rclone-wiz-application
  - id: linux-desktop-entry-exec
    desc: Desktop entry executable directive
    conf: 0.8
    crit: component
    any:
      - id: exec-eq-directive
      - id: hidden-false-directive
      - id: type-application-directive
    unless:
      - id: well-known/app/macos-benign::mattermost-desktop-bundle
      - id: well-known/app/macos-benign::mattermost-asar-resource
      - id: mattermost-desktop-resource-marker
      - id: well-known/app/rclone-wiz::rclone-wiz-application
  - id: linux-xdg-executable
    desc: XDG autostart with executable command
    crit: notable
    conf: 0.95
    mbc: F0006
    all:
      - id: linux-xdg-autostart-dir
      - id: linux-desktop-entry-exec
    unless:
      - id: well-known/app/macos-benign::mattermost-desktop-bundle
      - id: well-known/app/macos-benign::mattermost-asar-resource
      - id: mattermost-desktop-resource-marker
  - id: linux-multiple-persistence-mechanisms
    desc: Multiple Linux persistence mechanisms (XDG + Systemd)
    crit: hostile
    conf: 0.98
    attack: "T1547.013, T1543.002"
    all:
      - id: linux-xdg-executable
      - id: linux-user-systemd
