# objectives/persistence/system/surface/paths/linux.yaml
# Linux persistence surfaces that are not specific to cron, systemd services,
# SysV init, XDG autostart, or SSH authorized_keys.
defaults:
  platforms: [linux]
  for: [scripts, binaries]
  attack: T1547
traits:
- id: udev-rules-path
  desc: udev rules persistence path
  crit: component
  conf: 0.86
  attack: T1546
  if:
    type: text
    regex: '/(etc|lib|usr/lib)/udev/rules\.d/[^\s]{1,120}\.rules\b'
  unless:
  - type: text
    regex: '\$\((TARGET_DIR|DESTDIR)\)/(etc|lib|usr/lib)/udev/rules\.d/'

- id: udev-run-program
  desc: udev rule executes helper
  crit: notable
  conf: 0.86
  attack: T1546
  downgrade:
    any:
    - id: well-known/lib/xml2
    - id: well-known/lib/xslt
  if:
    type: text
    regex: '\b(RUN\+?=|PROGRAM=|IMPORT\{program\}=)'

- id: ld-so-conf-dropin-path
  desc: ld.so.conf.d linker path drop-in
  crit: notable
  conf: 0.84
  attack: T1574.006
  if:
    type: text
    regex: '/etc/ld\.so\.conf\.d/[^\s]{1,120}\.conf\b'

- id: polkit-rules-path
  desc: Polkit rules persistence path
  crit: suspicious
  conf: 0.86
  attack: T1548
  if:
    type: text
    regex: '/(etc|usr/share)/polkit-1/rules\.d/[^\s"]{1,120}\.rules([\s")''>;,]|$)'

- id: polkit-add-rule
  desc: Polkit addRule policy hook
  crit: suspicious
  conf: 0.88
  attack: T1548
  if:
    type: text
    substr: polkit.addRule

- id: dbus-service-path
  desc: D-Bus service activation path
  crit: suspicious
  conf: 0.84
  attack: T1546
  if:
    type: text
    regex: '/(etc|usr/share)/dbus-1/(system-)?services/[^\s]{1,120}\.service\b'

- id: dbus-system-policy-path
  desc: D-Bus system bus policy directory path
  crit: notable
  conf: 0.82
  attack: T1546
  # A path under the D-Bus system policy dir is a location signal, not intent:
  # this is the standard place packages install vendor .conf policy files.
  # Persistence/privesc is carried by an Exec= directive (dbus-exec-directive).
  if:
    type: text
    regex: '/(etc|usr/share)/dbus-1/system\.d/[^\s]{1,120}\.conf\b'

- id: dbus-exec-directive
  desc: D-Bus service Exec directive
  crit: component
  conf: 0.82
  attack: T1546
  if:
    type: text
    regex: '^\s*Exec=/[^\n]{1,200}'

- id: network-dispatcher-path--rx-1
  desc: Regex component marker
  crit: component
  conf: 0.86
  if:
    type: text
    regex: '/etc/(NetworkManager/dispatcher\.d|network/(if-up|if-down)\.d)/'

  count_max: 1032304
- id: network-dispatcher-path--rx-2
  desc: Regex component marker
  crit: component
  conf: 0.86
  if:
    type: text
    regex: '/etc/(dhcp/dhclient-(enter|exit)-hooks\.d)/'

  count_max: 1032305
- id: network-dispatcher-path--rx-3
  desc: Regex component marker
  crit: component
  conf: 0.86
  if:
    type: text
    regex: '/etc/(NetworkManager/dispatcher\.d|network/(if-pre-up|if-post-down)\.d)/'

  count_max: 1032306
- id: modules-load-path
  desc: Kernel module autoload path
  crit: notable
  conf: 0.84
  attack: T1547.006
  if:
    type: text
    regex: '/etc/modules-load\.d/[^\s]{1,120}\.conf\b|/etc/modules\b'
  unless:
  - id: well-known/tool/sysadmin/busybox::busybox-binary


- id: apt-hook-path
  desc: APT hook configuration path
  crit: notable
  conf: 0.84
  attack: T1546
  if:
    type: text
    regex: '/etc/apt/apt\.conf\.d/[^\s]{1,120}\b'

- id: apt-post-invoke-hook
  desc: APT post-invoke command hook
  crit: suspicious
  conf: 0.88
  attack: T1546
  if:
    type: text
    regex: 'APT::(Update::)?Post-Invoke|DPkg::(Pre|Post)-Invoke'

- id: user-mail-forward
  desc: User mail forward persistence file
  crit: notable
  conf: 0.82
  attack: T1546
  if:
    type: text
    regex: '(^|/)\.forward\b|/etc/aliases\b'
  unless:
  - id: micro-behaviors/dylib/library/family/

- id: supervisor-config-path
  desc: Supervisor service configuration path
  crit: notable
  conf: 0.84
  attack: T1543
  if:
    type: text
    regex: '/etc/(supervisor/conf\.d|supervisord\.d)/[^\s]{1,120}\.conf\b'
  unless:
    # AI-agent observability tools run their telemetry collector under supervisord.
  - id: well-known/tool/detection/ai-guardrail/

- id: runit-service-path
  desc: runit service persistence path
  crit: suspicious
  conf: 0.84
  attack: T1543
  if:
    type: text
    regex: '/etc/(sv|service)/[A-Za-z0-9_.@-]{1,80}/run\b'

- id: s6-service-path
  desc: s6 service persistence path
  crit: suspicious
  conf: 0.84
  attack: T1543
  if:
    type: text
    regex: '/etc/s6(-rc)?/[^\s]{1,160}/run\b'

- id: openrc-service-path
  desc: OpenRC service enablement
  crit: suspicious
  conf: 0.84
  attack: T1037.004
  if:
    type: text
    regex: '\brc-update\s+add\s+[A-Za-z0-9_.@-]{1,80}\b'
  unless:
    - type: text
      substr: VagrantPlugins
- id: managed-persistence-marker--rx-5
  desc: Regex component marker
  crit: component
  conf: 0.84
  if:
    type: text
    regex: '(?i)(managed|installed).{0,24}(service|cron)'

  count_max: 1032317
- id: managed-persistence-marker--rx-6
  desc: Regex component marker
  crit: component
  conf: 0.84
  if:
    type: text
    regex: '(?i)(managed|installed).{0,24}(init|autostart)'

  count_max: 1032318
  unless:
  - type: trait
    id: metadata/file/extension/identity::linux-kernel-package-artifact-basename
- id: managed-persistence-marker--rx-7
  desc: Regex component marker
  crit: component
  conf: 0.84
  if:
    type: text
    regex: '(?i)(persistence|autorun).{0,24}(service|cron)'

  count_max: 1032319
- id: managed-persistence-marker--rx-8
  desc: Regex component marker
  crit: component
  conf: 0.84
  if:
    type: text
    regex: '(?i)(persistence|autorun).{0,24}(init|autostart)'

  count_max: 1032320
- id: managed-persistence-marker--rx-9
  desc: Regex component marker
  crit: component
  conf: 0.84
  if:
    type: text
    regex: '(service|cron).{0,24}(managed|installed)'

  count_max: 1032321
- id: managed-persistence-marker--rx-10
  desc: Regex component marker
  crit: component
  conf: 0.84
  if:
    type: text
    regex: '(init|autostart).{0,24}(managed|installed)'

  count_max: 1032322
  unless:
  - type: trait
    id: metadata/file/extension/identity::linux-kernel-package-artifact-basename
- id: managed-persistence-marker--rx-11
  desc: Regex component marker
  crit: component
  conf: 0.84
  if:
    type: text
    regex: '(service|cron).{0,24}(persistence|autorun)'

  count_max: 1032323
- id: managed-persistence-marker--rx-12
  desc: Regex component marker
  crit: component
  conf: 0.84
  if:
    type: text
    regex: '(init|autostart).{0,24}(persistence|autorun)'

  count_max: 1032324
composite_rules:
- id: udev-exec-persistence
  desc: udev rule command persistence
  crit: suspicious
  conf: 0.9
  all:
  - id: udev-rules-path
  - id: udev-run-program

- id: pam-config-persistence
  desc: PAM configuration persistence
  crit: suspicious
  conf: 0.9
  all:
  - id: objectives/credential-access/pam/intercept::pam-d-write
  any:
  - id: objectives/credential-access/pam/intercept::pam-optional-inject
  - id: micro-behaviors/os/pam/auth::pam-optional-control

- id: polkit-rule-persistence
  desc: Polkit JavaScript rule persistence
  crit: suspicious
  conf: 0.9
  all:
  - id: polkit-rules-path
  - id: polkit-add-rule

- id: dbus-activation-persistence
  desc: D-Bus activation persistence
  crit: suspicious
  conf: 0.86
  all:
  - id: dbus-service-path
  - id: dbus-exec-directive

- id: linux-supervisor-persistence
  desc: Linux supervisor-managed persistence
  crit: suspicious
  conf: 0.86
  any:
  - id: supervisor-config-path
  - id: runit-service-path
  - id: s6-service-path
  - id: openrc-service-path
- id: stacked-linux-persistence-surfaces
  desc: Multiple Linux persistence surfaces
  crit: hostile
  conf: 0.96
  needs: 3
  any:
  - id: objectives/persistence/system/init/boot::linux-user-systemd
  - id: objectives/persistence/system/init/boot::linux-init-script
  - id: objectives/persistence/system/init/boot::linux-xdg-autostart
  - id: objectives/persistence/system/cron/job::crontab-reboot
  - id: objectives/persistence/login/shell/config::shell-files
  - id: objectives/evasion/kernel-hide/userspace::ld-so-preload-file
  - id: ld-so-conf-dropin-path
  - id: modules-load-path
  unless:
  - id: metadata/build/docker::dockerfile-basename
- id: managed-stacked-linux-persistence
  desc: Tagged stacked Linux persistence
  crit: hostile
  conf: 0.98
  all:
  - id: managed-persistence-marker
  - id: stacked-linux-persistence-surfaces
- id: network-dispatcher-path
  desc: Linux network dispatcher hook path
  crit: suspicious
  conf: 0.86
  attack: T1546
  any:
  - id: network-dispatcher-path--rx-1
  - id: network-dispatcher-path--rx-2
  - id: network-dispatcher-path--rx-3
  - id: objectives/persistence/system/surface/paths::network-dispatcher-path--rx-2
- id: managed-persistence-marker
  desc: Managed persistence artifact marker
  crit: notable
  conf: 0.84
  any:
  - id: managed-persistence-marker--rx-5
  - id: managed-persistence-marker--rx-6
  - id: managed-persistence-marker--rx-7
  - id: managed-persistence-marker--rx-8
  - id: managed-persistence-marker--rx-9
  - id: managed-persistence-marker--rx-10
  - id: managed-persistence-marker--rx-11
  - id: managed-persistence-marker--rx-12
