# objectives/command-and-control/dropper/execution/installer/7zip-sfx-bundle.yaml
defaults:
  platforms: [windows]
  for: [pe]
  attack: T1105
  mbc: B0022
traits:
  - id: large-overlay-bundle
    desc: Large overlay bundle component
    crit: component
    conf: 0.96
    size_min: 1000000
    if:
      id: metadata/binary/layout::large-overlay-1mb
    unless:
      - id: objectives/command-and-control/dropper/execution/installer/
    downgrade:
      any:
        - id: metadata/binary/installer/family/
        - id: metadata/binary/installer/script/
        - id: metadata/binary/installer/script/
        - id: metadata/binary/installer/sfx::7zsfx-config-marker
        - id: micro-behaviors/fs/temp/py2exe::py2exe-bundle
  - id: sfx-multiple-runprogram-stages
    desc: SFX launches multiple extracted executables
    crit: notable
    conf: 0.98
    if:
      type: text
      regex: 'RunProgram.{0,128}'
    count_min: 2
composite_rules:
  - id: admin-multi-stage-sfx-dropper
    desc: Admin SFX runs programs from packed overlay
    crit: hostile
    conf: 0.99
    attack: T1027.009,T1105,T1548.002
    all:
      - id: objectives/command-and-control/dropper/execution/installer::sfx-multiple-runprogram-stages
      - id: metadata/binary/layout::overlay-over-90pct
      - id: metadata/binary/layout::high-entropy-overlay
      - id: metadata/binary/layout::pe-checksum-zeroed
      - id: metadata/build/manifest::pe-requires-admin
      - id: micro-behaviors/process/create/shell/spawn::win-shellexecute
    unless:
      - id: metadata/binary/installer/family/
      - id: metadata/signed/platform/

  - id: signed-recognized-large-overlay-developer
    desc: Signed recognized large overlay developer
    crit: component
    conf: 0.95
    any:
      - id: metadata/signed/developer::win.rar-gmbh
      - id: metadata/signed/developer::valve-corp.
      - id: metadata/signed/developer::martin-tofall
      - id: metadata/signed/platform::microsoft
      - id: metadata/signed/platform::ip-zinukhov-stanislav-igorevich
  - id: large-overlay-installer-marker
    desc: Known installer marker for large overlay bundles
    crit: component
    conf: 0.9
    any:
      - id: metadata/binary/installer/family::known-benign-installer
      - id: metadata/binary/installer/script::inno-setup-installer
      - id: metadata/binary/installer/script::nsis
      - id: metadata/binary/installer/sfx::7zsfx-config-marker
      - id: micro-behaviors/fs/temp/py2exe::py2exe-bundle
  - id: signed-7zip-sfx-stage-loader
    desc: 7-Zip SFX bundle with staged loader behavior
    crit: suspicious
    conf: 0.98
    size_min: 1000000
    all:
      - id: metadata/binary/layout::large-overlay-1mb
      - id: metadata/binary/installer/sfx::7zip-sfx-runtime
      - id: micro-behaviors/fs/temp/directory::get-temp-path-dup
      - id: micro-behaviors/os/registry/access::registry-write-api-chain
    needs: 2
    any:
      - id: micro-behaviors/crypto/certificate/install::cert-install-api
      - id: micro-behaviors/os/sysinfo/platform/locale::locale-fingerprinting-imports
      - id: objectives/anti-analysis/debugger-detect/check::is-debugger-present-symbol
    unless:
      - id: metadata/signed/platform/
      - id: metadata/binary/installer/family/
  - id: signed-sfx-wrapper-layout
    desc: Signed SFX wrapper layout
    crit: suspicious
    conf: 0.97
    size_min: 900000
    size_max: 1300000
    all:
      - id: metadata/binary/layout::signed-pe-wrapper
      - id: metadata/binary/installer/sfx::7zip-sfx-archive-profile

  - id: signed-admin-pe-with-rsrc-payload
    desc: Admin-elevating PE bundling another PE in resources with timestamp anomaly
    crit: suspicious
    conf: 0.9
    all:
      - id: metadata/build/manifest::pe-requires-admin
      - id: micro-behaviors/data/embedded/payload/resource::pe-magic-in-rsrc
      - id: metadata/binary/anomaly/provenance::pe-far-future-timestamp
    unless:
      - id: metadata/signed/platform/
      - id: metadata/binary/installer/family/
  - id: identity-mismatched-7z-stage-loader
    desc: Identity-mismatched PE launches an appended 7-Zip stage
    crit: hostile
    conf: 0.99
    attack: T1036.005,T1027.009,T1105
    all:
      - id: objectives/evasion/masquerade/version-resource/contradiction::softwareok-free-manga-identity-contradiction
      - id: micro-behaviors/data/embedded/payload/archive::late-embedded-7z-archive
      - id: objectives/anti-static/obfuscation/payload/self-read::self-read-file-chain
      - id: micro-behaviors/process/create/shell/spawn::win-createprocess
    unless:
      - id: metadata/binary/installer/family::known-benign-installer
      - id: metadata/signed/platform/

  - id: unsigned-rar-elevated-logon-dropper
    desc: Unsigned PE embeds RAR payload and elevated logon task
    crit: hostile
    conf: 0.99
    attack: T1027.009,T1053.005,T1105
    all:
      - id: micro-behaviors/data/embedded/payload/archive::rar-magic
      - id: metadata/build/manifest::pe-requires-admin
      - id: objectives/persistence/login/scheduled-task/cli::windows-schtasks-onlogon-highest
      - id: metadata/signed/trust-level::unsigned
    unless:
      - id: metadata/signed/platform/
      - id: metadata/binary/installer/family/
      - id: metadata/binary/installer/script/
