# objectives/persistence/system/network-appliance/openwrt-python.yaml
defaults:
  platforms: [openwrt]
  for: [python, pyc]
  attack: T1543.002

traits:
  - id: openwrt-profile-path
    desc: OpenWrt shell profile path
    crit: notable
    conf: 0.84
    if:
      type: text
      substr: /etc/profile
    unless:
      - id: well-known/tool/detection/security-scanner::auditr-defensive-signature-catalog
      - id: metadata/package/quality/false-positives-scanner::installtime-audit-defensive-scanner-context
  - id: openwrt-opkg-postinst-path
    desc: OpenWrt opkg post-install hook path
    crit: notable
    conf: 0.88
    if:
      type: text
      substr: /etc/opkg/postinst.d
  - id: openwrt-hotplug-script-path
    desc: OpenWrt hotplug script path
    crit: notable
    conf: 0.88
    if:
      type: text
      regex: '/etc/hotplug\.d/[A-Za-z0-9_.-]+'
composite_rules:
  - id: openwrt-python-initd-service-persistence
    desc: OpenWrt Python init.d service persistence
    crit: hostile
    conf: 0.96
    all:
      - id: micro-behaviors/os/service/config::init-script-path
      - id: objectives/persistence/system/network-appliance::openwrt-rc-common-init
    any:
      - id: objectives/persistence/system/network-appliance::openwrt-initd-enable-command
      - id: micro-behaviors/process/create/subprocess::shell-true
      - id: micro-behaviors/process/create/subprocess::subprocess-shell-pyc
    needs: 2
  - id: openwrt-python-multi-surface-persistence
    desc: OpenWrt Python multi-surface persistence
    crit: hostile
    conf: 0.97
    any:
      - id: objectives/persistence/system/network-appliance::openwrt-root-crontab-path
      - id: objectives/persistence/system/init/boot::etc-rc-local
      - id: objectives/persistence/system/network-appliance::openwrt-uci-system-change
      - id: openwrt-profile-path
      - id: openwrt-opkg-postinst-path
      - id: objectives/persistence/system/network-appliance::openwrt-hotplug-iface-path
      - id: openwrt-hotplug-script-path
    needs: 3
  - id: openwrt-python-firmware-payload-persistence
    desc: OpenWrt Python firmware payload persistence
    crit: hostile
    conf: 0.96
    all:
      - id: micro-behaviors/os/service/config::init-script-path
      - id: objectives/persistence/system/network-appliance::openwrt-rc-common-init
    any:
      - id: objectives/persistence/system/init/boot::etc-rc-local
      - id: objectives/persistence/system/network-appliance::openwrt-uci-system-change
      - id: objectives/persistence/system/network-appliance::openwrt-hotplug-iface-path
      - id: openwrt-hotplug-script-path
    needs: 2
