# objectives/command-and-control/backdoor/script/powershell-rat.yaml
# PowerShell RAT — Reflective Loader + Registry Persistence + CIM Recon
# Detects the multi-platform RAT pattern delivered as a dropped PowerShell script:
# HKCU\Run persistence via BAT, in-memory .NET assembly loading via Reflection,
# hidden child PowerShell for sub-command execution, and CIM process enumeration.
# ATT&CK: T1059.001 (PowerShell), T1620 (Reflective Code Loading), T1547.001 (Run Key)
defaults:
  platforms:
    - windows
  for:
    - powershell
  size_max: 5242880
  attack: T1059.001
traits:
  - id: ps1-assembly-load-b64
    desc: In-memory .NET assembly load from base64
    crit: suspicious
    conf: 0.95
    attack: T1620
    if:
      type: text
      regex: \[System\.Reflection\.Assembly\]::Load.{0,80}FromBase64String
  - id: ps1-cim-process-enum
    desc: CIM Win32_Process enumeration for process list
    crit: notable
    conf: 0.85
    attack: T1057
    if:
      type: text
      regex: (?i)(Get-CimInstance|Get-WmiObject).{0,80}Win32_Process.{0,80}ForEach-Object
  - id: ps1-hidden-child-ps
    desc: Spawns hidden child PowerShell with redirected output
    crit: suspicious
    conf: 0.9
    attack: T1564.003
    if:
      type: text
      regex: ProcessWindowStyle\]::Hidden.{0,120}RedirectStandardOutput
  - id: ps1-webclient-upload
    desc: WebClient UploadData for C2 beacon
    crit: notable
    conf: 0.8
    attack: T1071.001
    if:
      type: text
      substr: UploadData
  - id: ps1-c2-api-poll
    desc: C2 client poll API endpoint pattern
    crit: notable
    conf: 0.88
    attack: T1071.001
    if:
      type: text
      substr: /api/client/poll/
  - id: ps1-c2-api-response
    desc: C2 client response API endpoint pattern
    crit: component
    conf: 0.85
    if:
      type: text
      substr: /api/client/response
  - id: ps1-default-servers
    desc: defaultServers C2 server list variable
    crit: component
    conf: 0.85
    if:
      type: text
      word: defaultServers
  - id: ps1-temp-guid-ps1
    desc: Temp GUID ps1 path
    crit: suspicious
    conf: 0.9
    if:
      type: text
      regex: '\$env:TEMP\\c\$\(\[guid\]::NewGuid\(\)\)\.ps1'
  - id: ps1-temp-guid-txt
    desc: Temp GUID txt path
    crit: suspicious
    conf: 0.88
    if:
      type: text
      regex: '\$env:TEMP\\o\$\(\[guid\]::NewGuid\(\)\)\.txt'
  - id: ps1-temp-hidden-startprocess
    desc: Hidden temp PowerShell runner
    crit: suspicious
    conf: 0.93
    if:
      type: text
      regex: 'Start-Process powershell.{0,200}Out-File.{0,120}-WindowStyle Hidden'
  - id: ps1-waitforexit-timeout
    desc: WaitForExit timeout call
    crit: suspicious
    conf: 0.88
    if:
      type: text
      substr: WaitForExit($config.commandTimeout * 1000)
  - id: ps1-stop-process-force-id
    desc: Stop-Process force by id
    crit: suspicious
    conf: 0.88
    if:
      type: text
      regex: 'Stop-Process\s+-Id\s+\$[A-Za-z_][A-Za-z0-9_.]*\s+-Force'
  - id: ps1-parentprocessid-filter
    desc: ParentProcessId kill filter
    crit: suspicious
    conf: 0.9
    if:
      type: text
      regex: 'ParentProcessId=\$\(\$[A-Za-z_][A-Za-z0-9_.]*\)'
  - id: ps1-schtask-action-wscript
    desc: Scheduled task runs wscript
    crit: suspicious
    conf: 0.92
    if:
      type: text
      regex: 'New-ScheduledTaskAction\s+-Execute\s+"wscript\.exe"'
  - id: ps1-register-schtask
    desc: Registers scheduled task
    crit: suspicious
    conf: 0.88
    unless:
      - id: well-known/tool/sysadmin/kubernetes::minikube
    if:
      type: text
      regex: '(?:^|[^\{A-Za-z0-9_-])Register-ScheduledTask\b'
      case_insensitive: true
  - id: ps1-vbs-update-launch
    desc: VBS update.ps1 launcher
    crit: suspicious
    conf: 0.9
    if:
      type: text
      regex: 'CreateObject\("WScript\.Shell"\).{0,200}update\.ps1'
  - id: ps1-vbs-ps1-fanout
    desc: VBS iterates sibling ps1
    crit: suspicious
    conf: 0.92
    if:
      type: text
      regex: 'For Each .{0,160}GetExtensionName\(.{0,80}\)\) = "ps1".{0,160}objShell\.Run'
  - id: ps1-status-beacon-uri
    desc: Launch-status beacon URI path
    crit: suspicious
    conf: 0.95
    if:
      type: text
      regex: '(?i)https?://[0-9a-z.\[\]-]{3,120}/stuk-phase\b'
  - id: ps1-status-beacon-message
    desc: Launch-status beacon message family
    crit: suspicious
    conf: 0.97
    count_min: 2
    if:
      type: text
      regex: '\b(?:GFILE FOUND ON PC|RDOWNLOAD ERROR|RFATAL DOWNLOAD ERROR|GLAUNCH SUCCESS|RLAUNCH FAILED|GSESSION CLOSED)\b'
  - id: ps1-status-beacon-post
    desc: Invoke-RestMethod status POST body
    crit: suspicious
    conf: 0.94
    if:
      type: text
      regex: '(?is)Invoke-RestMethod.{0,160}-Method\s+Post.{0,160}-Body\s+@\{.{0,160}\bmessage\s*=.{0,120};.{0,120}\btag\s*='
  - id: ps1-started-launch-branch
    desc: Started branch posts launch result
    crit: suspicious
    conf: 0.88
    if:
      type: text
      regex: '(?is)if\s*\(\s*\$started\s*\).{0,240}GLAUNCH SUCCESS.{0,240}else.{0,240}RLAUNCH FAILED'
composite_rules:
  - id: ps1-aes-beacon-rat
    desc: PowerShell RAT with AES-encrypted payload and C2 beacon polling
    crit: hostile
    conf: 0.96
    attack: T1071.001
    all:
      - id: objectives/command-and-control/dropper/staging/encrypted::ps-aes-inline-payload
    any:
      - id: ps1-c2-api-poll
      - id: ps1-c2-api-response
      - id: ps1-default-servers
  - id: ps1-reflective-c2-agent
    desc: PowerShell RAT with reflective .NET load and registry persistence
    crit: hostile
    conf: 0.97
    all:
      - id: ps1-assembly-load-b64
      - id: ps1-cim-process-enum
      - id: ps1-hidden-child-ps
      - id: objectives/persistence/login/registry/autostart::run-key
    any:
      - id: ps1-webclient-upload
    needs: 1
  - id: ps1-polling-command-runner
    desc: PowerShell polling backdoor with hidden temp-script command runner
    crit: hostile
    conf: 0.98
    attack: T1071.001
    near_lines: 40
    all:
      - id: ps1-temp-guid-ps1
      - id: ps1-temp-guid-txt
      - id: ps1-temp-hidden-startprocess
      - id: ps1-waitforexit-timeout
      - id: ps1-stop-process-force-id
      - id: ps1-parentprocessid-filter
    any:
      - id: ps1-c2-api-poll
      - id: ps1-c2-api-response
      - id: micro-behaviors/communications/http/lib/invoke-webrequest::ps-post
    needs: 1
  - id: ps1-wscript-schtask-backdoor
    desc: PowerShell backdoor persists via scheduled wscript launcher
    crit: hostile
    conf: 0.98
    attack: T1053.005
    near_lines: 80
    all:
      - id: ps1-schtask-action-wscript
      - id: ps1-register-schtask
      - id: ps1-vbs-update-launch
      - id: ps1-vbs-ps1-fanout
    any:
      - id: objectives/execution/interpreter/vbscript/interpreter::vbs-powershell-launcher
      - id: micro-behaviors/process/create/shell/wsh::wscript-exe
      - id: micro-behaviors/os/com/object::wscript-shell-literal
    needs: 1
  - id: ps1-launch-status-beacon
    desc: PowerShell loader posts launch status beacons
    crit: hostile
    conf: 0.98
    attack: T1071.001
    near_lines: 40
    all:
      - id: ps1-status-beacon-uri
      - id: ps1-status-beacon-message
      - id: ps1-status-beacon-post
    any:
      - id: ps1-started-launch-branch
      - id: micro-behaviors/communications/http/lib/invoke-webrequest::ps-post
      - id: micro-behaviors/process/create::start-process-dup
    needs: 1
