# objectives/persistence/login/startup/registry/windows.yaml
defaults:
  attack: T1547.001
  crit: suspicious
  for:
  - pe
  - oledoc
  - powershell
  - batch
  - csharp
  - python
  platforms:
  - windows
traits:
- id: start-menu-startup
  desc: Start Menu Programs Startup path
  crit: notable
  conf: 0.7
  if:
    type: text
    regex: \\{1,2}Start Menu\\{1,2}Programs\\{1,2}Startup
  downgrade:
    any:
    - id: metadata/binary/vendor::microsoft-corp-versioninfo
  unless:
  - id: objectives/command-and-control/dropper/execution/resource/
- id: appdata-startup
  desc: AppData Roaming Startup path
  crit: notable
  conf: 0.7
  if:
    type: text
    regex: Roaming.{1,50}Programs.{1,20}Startup
- id: shell-startup
  desc: shell:startup special folder
  crit: notable
  conf: 0.7
  if:
    type: text
    substr: shell:startup
- id: software-run-key
  desc: References SOFTWARE Run path only
  crit: component
  conf: 0.7
  if:
    type: text
    regex: (?i)SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run
  unless:
  - id: metadata/signed/developer::tencent-technologyshenzhen-company-limited
  - id: well-known/lib/wmi::wmi-win32-provider-versioninfo
  - id: metadata/binary/installer/family/
  - id: objectives/command-and-control/dropper/execution/resource/
- id: hkcu-run
  desc: References HKCU Run path only
  crit: component
  conf: 0.6
  if:
    type: text
    regex: HKCU.{0,50}\\Run(\\|$|["'\s])
  unless:
  - id: well-known/tool/detection/malware-config-parser
- id: hklm-run
  desc: References HKLM Run path only
  crit: component
  conf: 0.6
  if:
    type: text
    regex: HKLM.{0,50}\\Run(\\|$|["'\s])
  unless:
  - id: well-known/tool/detection/malware-config-parser
  - id: metadata/package/manifest/vendor::bitdefender-ctc-module
  - id: metadata/lang/source/
  - id: metadata/lang/source/

- id: script-style-run-value
  desc: Script-like Run value name
  crit: component
  conf: 0.78
  if:
    type: text
    regex: (?i)\b[a-z]{4,}[_-]script\d{2,}\b
- id: run-key-wscript-vbs-launcher
  desc: Run key launches WScript VBS
  crit: suspicious
  conf: 0.95
  if:
    type: text
    regex: (?is)CurrentVersion\\Run.{0,120}wscript\.exe.{0,220}\.vbs
- id: runonce-word
  desc: References RunOnce registry path only
  crit: component
  conf: 0.5
  if:
    type: text
    regex: (?i)(\\|/)RunOnce(\\|$)
  unless:
  - id: objectives/evasion/self-delete/cleanup/heartbeat::advanced-inf-setup-regpath
  - id: metadata/signed/platform::microsoft-corporation
  - id: metadata/signed/leaf::microsoft-corporation
  - id: metadata/binary/installer/family/
  - id: metadata/binary/vendor::python-software-foundation-marker
  - id: metadata/signed/developer::python-software-foundation
  - id: metadata/signed/leaf::python-software-foundation
- id: import-winreg
  desc: import winreg statement
  crit: component
  conf: 0.5
  if:
    type: text
    substr: import winreg
- id: from-winreg
  desc: from winreg import
  crit: component
  conf: 0.5
  if:
    type: text
    substr: from winreg
- id: underscore-winreg
  desc: _winreg module (Python 2)
  crit: component
  conf: 0.5
  if:
    type: text
    substr: _winreg
- id: openkey-hkey
  desc: OpenKey registry path lookup
  crit: component
  conf: 0.6
  if:
    type: text
    regex: \bOpenKey(Ex)?\s*\(\s*(winreg\.)?HKEY_[A-Z_]+[\s\S]{0,220}
- id: setvalueex
  desc: SetValueEx call
  crit: component
  conf: 0.6
  for:
  - powershell
  - batch
  if:
    type: string_literal
    regex: SetValueEx\b
- id: createkey
  desc: CreateKey call
  crit: component
  conf: 0.6
  for:
  - powershell
  - batch
  if:
    type: string_literal
    substr: CreateKey
- id: installservice-fn
  desc: InstallService function
  crit: notable
  conf: 0.7
  if:
    type: text
    regex: InstallService\s*\(
- id: windows-powershell-registry
  desc: Modifies registry for persistence via PowerShell
  conf: 0.9
  crit: suspicious
  if:
    type: trait
    id: objectives/persistence/login/registry/autostart::run-key-write-ps
- id: game-service-run-value-prefix
  desc: Game service Run value prefix
  crit: component
  conf: 0.94
  if:
    type: text
    word: ServizioGioco_
- id: network-module-run-value-prefix
  desc: Network module Run value prefix
  crit: component
  conf: 0.94
  if:
    type: text
    word: ModuloRete_
- id: italian-run-key-success
  desc: Italian Run key success text
  crit: component
  conf: 0.9
  if:
    type: text
    substr: aggiunto all'avvio
- id: six-py-basename
  desc: six.py compatibility file
  crit: notable
  for: [python]
  if:
    type: path
    basename: true
    exact: six.py
composite_rules:
- id: createservice-api
  desc: CreateService install flow
  crit: suspicious
  conf: 0.8
  all:
  - id: micro-behaviors/os/service/control::create-service
  - id: micro-behaviors/os/service/control::open-sc-manager
  any:
  - id: micro-behaviors/os/service/control::sc-create-command
  - id: installservice-fn
  - id: micro-behaviors/os/service/control::python-win32serviceutil-module
- id: windows-startup-folder
  desc: Targets Windows Startup folder for persistence
  conf: 0.85
  crit: notable
  any:
  - id: start-menu-startup
  - id: appdata-startup
  - id: shell-startup
- id: windows-registry-run
  desc: Accesses Windows registry Run key
  conf: 0.9
  crit: notable
  any:
  - id: software-run-key
  - id: hkcu-run
- id: windows-run-key-script-launcher
  desc: Run key persists script launcher
  crit: suspicious
  conf: 0.96
  all:
  - id: windows-registry-run
  - id: run-key-wscript-vbs-launcher
  - id: hklm-run
  - id: micro-behaviors/os/registry/keys::currentversion-run
  - id: micro-behaviors/os/registry/keys::windows-registry-run-access
  any:
  - id: script-style-run-value
  downgrade:
    any:
    - id: metadata/package/quality/profile/
- id: windows-registry-windows-run-value
  desc: Writes Run value in Windows key
  conf: 0.9
  crit: suspicious
  all:
  - id: micro-behaviors/os/registry/keys::currentversion-windows-load-run-path
  - id: micro-behaviors/os/registry/keys::run-value-name
  - id: micro-behaviors/os/registry/access::reg-set-value-ex-w
- id: windows-registry-runonce
  desc: Accesses Windows registry RunOnce key
  conf: 0.9
  crit: notable
  any:
  - id: objectives/persistence/login/registry/autostart::run-once-key
  - id: runonce-word
  unless:
  - id: metadata/binary/installer/family/
- id: windows-python-winreg-import
  desc: Imports Python winreg module for
  conf: 0.7
  crit: suspicious
  any:
  - id: import-winreg
  - id: from-winreg
  - id: underscore-winreg
  unless:
  - id: well-known/app/turboadb::official-turboadb-gui
- id: windows-python-winreg
  desc: Accesses Windows registry via Python
  conf: 0.8
  crit: suspicious
  size_max: 5000
  any:
  - id: import-winreg
  - id: from-winreg
  - id: underscore-winreg
  - id: openkey-hkey
  - id: setvalueex
  - id: createkey
  unless:
  - id: six-py-basename
  - id: well-known/app/turboadb::official-turboadb-gui
- id: windows-service-install
  desc: Installs Windows service for persistence
  conf: 0.85
  attack: T1543.003
  any:
  - id: micro-behaviors/os/service/control::sc-create-command
  - id: createservice-api
  - id: installservice-fn
- id: service-persistence-dropper-loader
  desc: Service persistence paired with dropper loader
  crit: hostile
  conf: 0.96
  for: [pe]
  all:
  - id: micro-behaviors/os/service/control::sc-create-command
  - id: objectives/command-and-control/dropper/download::winhttp-temp-dropper
  - id: objectives/evasion/security-bypass/policy/execution::bypass-execution-policy
  any:
  - id: metadata/build/manifest::pe-requires-admin
  - id: objectives/command-and-control/beacon/network/http::direct-http-c2-binary-ip
  - id: objectives/anti-analysis/emulator-detect/vendor-strings::comprehensive
  downgrade:
    any:
    - id: metadata/signed/platform/
    - id: metadata/binary/installer/family/
- id: windows-wmi-persistence
  desc: Sets up WMI event subscription
  conf: 0.9
  crit: suspicious
  attack: T1546.003
  any:
  - id: objectives/persistence/system/wmi/subscription::set-wmiinstance-eventfilter
  - id: objectives/persistence/system/wmi/subscription::set-wmiinstance-commandline-consumer
  - id: objectives/persistence/system/wmi/subscription::set-wmiinstance-filter-binding
- id: pe-startup-folder-registry-persistence
  desc: Startup folder persistence with registry fingerprinting
  crit: suspicious
  conf: 0.9
  all:
  - id: objectives/persistence/login/startup/folder::startup-folder-shell-path-composition-chain
  - id: micro-behaviors/os/registry/access::registry-write-api-chain
  any:
  - id: micro-behaviors/process/sync/mutex::guid-mutex
  - id: micro-behaviors/process/info/current::get-module-filename
  - id: micro-behaviors/process/create/spawn::create-process-w
  downgrade:
    any:
    - id: metadata/signed/platform/
    - id: metadata/binary/vendor::is-microsoft-named
    - id: metadata/binary/metrics/structural::many-imports
    - id: micro-behaviors/os/service/config/
    - id: metadata/package/versioning::pe-version-resource
    - id: metadata/lang/compiler/native/
- id: pe-run-key-self-install
  desc: PE self-copies and persists via Run key
  crit: suspicious
  conf: 0.92
  for: [pe, dll]
  all:
  - id: micro-behaviors/os/registry/access::registry-write-api-chain
  - id: software-run-key
  any:
  - id: micro-behaviors/fs/file/copy::copyfile-a
  - id: micro-behaviors/os/sysinfo/directories::get-system-directory
  - id: micro-behaviors/fs/path/check::get-module-filename
  needs: 2
  downgrade:
    any:
    - id: metadata/signed/platform/
    - id: metadata/binary/vendor::is-microsoft-named
    - id: metadata/binary/metrics/structural::many-imports
    - id: micro-behaviors/os/service/config/
- id: pe-wide-self-copy-run-key-install
  desc: PE self-copy with Run-key persistence
  crit: hostile
  conf: 0.97
  mbc: F0006
  for: [pe, dll]
  all:
  - id: software-run-key
  - id: micro-behaviors/os/registry/access::registry-write-api-chain
  - id: micro-behaviors/fs/file/copy::copyfile-w
  - id: micro-behaviors/fs/path/check::get-module-filename
  unless:
  - id: metadata/signed/platform/
  - id: metadata/binary/installer/family/
  - id: metadata/package/versioning::pe-version-resource
- id: pe-game-service-run-key-persistence
  desc: Game service Run key persistence
  crit: suspicious
  conf: 0.96
  for: [pe]
  all:
  - id: micro-behaviors/os/registry/keys::windows-registry-run-access
  - id: micro-behaviors/os/registry/manipulate::registry-write-api-chain
  any:
  - id: game-service-run-value-prefix
  - id: network-module-run-value-prefix
  - id: italian-run-key-success
  needs: 2
  unless:
  - id: metadata/signed/platform/
  - id: metadata/binary/vendor::microsoft-corp-versioninfo
- id: raw-ip-run-key-stager
  desc: Run-key persistence for raw-IP stager
  crit: hostile
  conf: 0.97
  ref: sha256:8a710cccc278dc81e3fda6081c71cfac7bcfa0099821d7d0a7fab7f415cfd411
  for: [pe, dll]
  all:
  - id: software-run-key
  - id: micro-behaviors/os/registry/access::registry-write-api-chain
  - id: micro-behaviors/communications/http/url/external-ip::url-external-ip-port-binary
  any:
  - id: objectives/command-and-control/dropper/delivery/download-execute::dynamic-winhttp-downloader
  - id: objectives/command-and-control/dropper/delivery/download-execute::raw-ip-winhttp-memory-stager
  - id: micro-behaviors/process/create/spawn::create-process-a
  - id: micro-behaviors/process/create/spawn::create-process-w
  needs: 2
  downgrade:
    any:
    - id: metadata/signed/platform/
    - id: metadata/binary/installer/family/
- id: vb6-run-key-self-install
  desc: VB6 self-copy Run key persistence
  crit: hostile
  conf: 0.96
  for: [pe]
  all:
  - id: metadata/binary/framework::vb6-application-framework
  - id: micro-behaviors/os/registry/manipulate::registry-write-api-chain
  - id: objectives/persistence/login/registry/autostart::run-key
  - id: micro-behaviors/fs/file/copy::vb6-rtcfilecopy
  - id: micro-behaviors/os/sysinfo/directories::get-windows-directory
  - id: micro-behaviors/os/sysinfo/directories::get-system-directory-string
  unless:
  - id: metadata/signed/platform/
  - id: metadata/binary/vendor::is-microsoft-named
