# objectives/persistence/system/launchd/source/unix-source.yaml
defaults:
  platforms: [unix]
  for: [go, kotlin, class, jar, php, ruby, rust, javascript]
  attack: T1543.004

traits:
  - id: source-launchagents-join
    desc: Source LaunchAgents join
    crit: notable
    conf: 0.92
    if:
      type: text
      regex: '["'']Library["''][,\)][\s\S]{0,80}["'']LaunchAgents["'']'

  - id: source-launchagents-plist
    desc: Source LaunchAgents plist
    crit: notable
    conf: 0.92
    if:
      type: text
      regex: 'Library/LaunchAgents/.{0,80}\.plist'

  - id: source-launchctl-load
    desc: Source launchctl load call
    crit: notable
    conf: 0.92
    if:
      type: text
      regex: '\blaunchctl\b.{0,80}\bload\b'
    unless:
      - type: text
        substr: com.apple.smb.preferences

  - id: source-runatload-key
    desc: Source RunAtLoad enabled value
    crit: notable
    conf: 0.9
    if:
      type: text
      regex: 'RunAtLoad[\s\S]{0,80}(<true\s*/>|true\b)'

composite_rules:
  - id: source-launchd-user-persistence
    desc: Source installs enabled LaunchAgent
    crit: hostile
    conf: 0.97
    all:
      - id: source-launchctl-load
      - id: source-runatload-key
    any:
      - id: source-launchagents-join
      - id: source-launchagents-plist
    needs: 1
    unless:
      - id: metadata/package/quality/false-positives::benign-build-tooling-context
      - id: well-known/tool/sysadmin/alpine-packages::alpine-homebrew-package
      - id: well-known/lib/crate-package/
