# objectives/persistence/system/launchd/source/unix-source.yaml
defaults:
  platforms: [unix]
  for: [go, kotlin, php, ruby, rust, javascript, python]
  attack: T1543.004

composite_rules:
  - id: source-launchd-user-persistence
    desc: Source loads user LaunchAgent
    # User LaunchAgents are a persistence mechanism but are also routinely
    # installed by legitimate desktop applications and daemon managers.
    crit: notable
    conf: 0.97
    all:
      - id: micro-behaviors/os/service/control::launchctl-load-command
    any:
      - id: micro-behaviors/os/service/config::source-launchagents-join
      - id: micro-behaviors/os/service/config::source-launchagents-plist
    needs: 1
    unless:
      - id: metadata/package/quality/false-positives::benign-build-tooling-context
      - id: well-known/tool/sysadmin/alpine-packages::alpine-homebrew-package
      - id: well-known/lib/crate-package/
      - id: well-known/lib/running-process::runpm-launchd-backend
      - id: well-known/app/spacedrive/
      - id: well-known/app/pixcode::benign-runtime-context
      - id: well-known/app/pixcode::benign-archive-context
      - id: well-known/app/pixcode::daemon-manager-source
      - id: well-known/app/pixcode::archive-member-context
  - id: launchd-agent-persistence-with-claude-config
    desc: LaunchAgent persistence paired with Claude Code configuration writes
    # This combination is security-relevant but has a documented legitimate
    # use in AI proxy software that installs a user service and points the
    # client at its local endpoint. Keep it suspicious outside a recognized
    # package context rather than treating the capability combination as an
    # attack with no legitimate use.
    crit: suspicious
    conf: 0.99
    all:
      - id: micro-behaviors/os/service/control::launchctl-load-command
      - id: micro-behaviors/os/service/config::source-launchagents-join
      - id: micro-behaviors/os/service/config::source-runatload-key
      - id: micro-behaviors/fs/path/application/config::claude-config-directory
      - id: micro-behaviors/fs/file/write/direct::node-write-file-sync-call
    unless:
      - id: well-known/app/neurolink::benign-archive-context
      - id: well-known/app/neurolink::proxy-launchd-context
      - id: well-known/app/agentgem::benign-archive-context
      - id: well-known/app/agentgem::warm-service-label
  - id: source-ssh-backdoor-launchd-persistence
    desc: SSH backdoor paired with launchd enablement
    crit: hostile
    conf: 0.99
    attack: "T1098.004,T1543.004"
    all:
      - id: objectives/persistence/login/ssh/backdoor::js-authorized-key-literal-backdoor
      - id: micro-behaviors/os/service/control::launchctl-load-command
      - id: objectives/persistence/system/surface/paths::remote-login-enabled
  - id: disguised-recurring-launchd-persistence
    desc: Disguised recurring LaunchAgent persistence
    crit: hostile
    conf: 0.98
    all:
      - id: micro-behaviors/os/service/control::launchctl-load-command
      - id: micro-behaviors/os/service/config::source-launchagents-join
      - id: micro-behaviors/os/service/config::source-runatload-key
      - id: micro-behaviors/os/service/config::source-startinterval-key
      - id: objectives/evasion/masquerade/binary::apple-helper-name
