# objectives/persistence/system/launchd/source/unix-source.yaml
defaults:
  platforms: [unix]
  for: [go, kotlin, php, ruby, rust, javascript, python]
  attack: T1543.004

traits:
  - id: source-launchagents-join
    desc: Source LaunchAgents join
    crit: notable
    conf: 0.92
    if:
      type: text
      regex: '["'']Library["''][\s\S]{0,80}["'']LaunchAgents["'']'

  - id: source-launchagents-plist
    desc: Source LaunchAgents plist
    crit: notable
    conf: 0.92
    if:
      type: text
      regex: 'Library/LaunchAgents/.{0,80}\.plist'

  - id: source-runatload-key
    desc: Source RunAtLoad enabled value
    crit: notable
    conf: 0.9
    if:
      type: text
      regex: 'RunAtLoad[\s\S]{0,80}(<true\s*/>|true\b)'

  - id: source-startinterval-key
    desc: Source recurring LaunchAgent interval
    crit: notable
    conf: 0.9
    if:
      type: text
      word: StartInterval

composite_rules:
  - id: source-launchd-user-persistence
    desc: Source loads user LaunchAgent
    # User LaunchAgents are a persistence mechanism but are also routinely
    # installed by legitimate desktop applications and daemon managers.
    crit: notable
    conf: 0.97
    all:
      - id: micro-behaviors/os/service/control::launchctl-load-command
    any:
      - id: source-launchagents-join
      - id: source-launchagents-plist
    needs: 1
    unless:
      - id: metadata/package/quality/false-positives::benign-build-tooling-context
      - id: well-known/tool/sysadmin/alpine-packages::alpine-homebrew-package
      - id: well-known/lib/crate-package/
      - id: well-known/lib/running-process::runpm-launchd-backend
      - id: well-known/app/spacedrive/
      - id: well-known/app/pixcode::benign-runtime-context
      - id: well-known/app/pixcode::benign-archive-context
      - id: well-known/app/pixcode::daemon-manager-source
      - id: well-known/app/pixcode::archive-member-context
  - id: source-ssh-backdoor-launchd-persistence
    desc: SSH backdoor paired with launchd enablement
    crit: hostile
    conf: 0.99
    attack: "T1098.004,T1543.004"
    all:
      - id: objectives/persistence/login/ssh/backdoor::js-authorized-key-literal-backdoor
      - id: micro-behaviors/os/service/control::launchctl-load-command
      - id: objectives/persistence/system/surface/paths::remote-login-enabled
  - id: disguised-recurring-launchd-persistence
    desc: Disguised recurring LaunchAgent persistence
    crit: hostile
    conf: 0.98
    all:
      - id: micro-behaviors/os/service/control::launchctl-load-command
      - id: source-launchagents-join
      - id: source-runatload-key
      - id: source-startinterval-key
      - id: objectives/evasion/masquerade/binary::apple-helper-name
