# micro-behaviors/os/autorun/scheduled/persistence-markers.yaml
# Persistence marker strings and Windows Maintenance path references.
# These are attack-context-specific fragments that have no meaning outside
# of malware persistence mechanisms. Combined in objectives/persistence/
# composites for full detection.
defaults:
  platforms: [windows]
  for: [pe, dll]
  crit: notable
  conf: 0.9
  attack: T1053.005
  mbc: F0012
traits:
  - id: persist-marker
    desc: Persistence configuration marker string
    if:
      type: text
      exact: "[PERSIST:1]"
  - id: windows-maintenance-path
    desc: Windows Maintenance scheduled task path reference
    if:
      type: text
      substr: "\\Microsoft\\Windows\\Maintenance\\"
