# objectives/persistence/system/service/systemd/obfuscated-bash-builder.yaml
# Bash systemd-user persistence builder embedded in obfuscated JS/script.
# JS-string-array obfuscators split a multi-line bash script into many
# small fragments, each rotated through a decoder. The reassembled
# script writes a [Service] unit to ~/.config/systemd/user/, sets
# Restart=always + WantedBy=default.target, and `systemctl --user
# enable && start`s it. The individual fragments here are uncommon
# in legitimate code; the cluster is hostile.
#
# Examples (npm xmlrpc, xmrdropper variants):
#   `\nSYSTEMD_PATH=$HOME/.config/systemd/user`
#   `Restart=always\nWantedBy=default.target`
#   `systemctl --user enable $APPNAME.service`
defaults:
  platforms: [unix]
  for: [scripts]
  attack: T1543.002
  mbc: F0011
traits:
- id: obfuscated-systemd-user-path-fragment
  desc: Obfuscated systemd user path fragment
  crit: component
  conf: 0.92
  # `systemd/us[^e]` matches the path only when it is broken before `user`
  # (a split-string concatenation fragment such as "systemd/us"+"er"), so the
  # ordinary contiguous `/usr/lib/systemd/user/` path is not flagged as
  # "obfuscated". The `\x0aSYSTEMD` leg is the literal newline-escape tell.
  if:
    type: text
    regex: '\\x0aSYSTEMD|systemd/us[^e]'
- id: systemctl-user-fragment--rx-1
  desc: systemctl --user command fragment
  crit: component
  conf: 0.92
  if:
    type: text
    regex: '\bsystemctl\s+--user\b'
  count_max: 1029793
- id: systemctl-user-fragment--rx-2
  desc: Regex component marker
  crit: component
  conf: 0.92
  if:
    type: text
    substr: systemctl\\x20--u
  count_max: 1029794
- id: systemctl-user-fragment--rx-3
  desc: Regex component marker
  crit: component
  conf: 0.92
  if:
    type: text
    substr: --user\\x20enab
  count_max: 1029795
- id: systemctl-user-fragment--rx-4
  desc: Regex component marker
  crit: component
  conf: 0.92
  if:
    type: text
    substr: user\\x20stop
  count_max: 1029796
- id: systemctl-user-fragment--rx-5
  desc: Regex component marker
  crit: component
  conf: 0.92
  if:
    type: text
    substr: user\\x20start
  count_max: 1029797
- id: service-unit-write-fragment
  desc: Service unit write redirection fragment
  crit: suspicious
  conf: 0.9
  if:
    type: text
    regex: '\.service\\\\x20>|E\.service\\\\x20'
- id: restart-always-fragment
  desc: Restart=always systemd directive fragment
  crit: suspicious
  conf: 0.92
  if:
    type: text
    regex: '\b(start|rt)=alwa'
- id: wanted-by-default-target-fragment
  desc: WantedBy default.target fragment
  crit: component
  conf: 0.92
  if:
    type: text
    regex: 'WantedBy=default\.target|t\.target\\\\x0a|=defaul.{0,5}target'
- id: bash-os-release-read
  desc: Reads /etc/os-release for arch detection
  crit: component
  conf: 0.85
  if:
    type: text
    regex: '/etc/os-re|/etc/os-release|os-release'
composite_rules:
- id: obfuscated-systemd-bash-builder
  desc: Embedded systemd persistence bash builder
  crit: hostile
  conf: 0.96
  all:
  - id: systemctl-user-fragment
  - id: service-unit-write-fragment
  needs: 2
  any:
  - id: obfuscated-systemd-user-path-fragment
  - id: micro-behaviors/os/service/config::systemd-user-service-directory-path
  - id: restart-always-fragment
  - id: objectives/persistence/system/service/systemd::restart-always-text
  - id: wanted-by-default-target-fragment
  - id: bash-os-release-read
  downgrade:
    any:
      - id: micro-behaviors/os/service/legitimate::installer-script
- id: systemctl-user-fragment
  desc: systemctl --user command fragment
  crit: component
  conf: 0.92
  any:
  - id: systemctl-user-fragment--rx-1
  - id: systemctl-user-fragment--rx-2
  - id: systemctl-user-fragment--rx-3
  - id: systemctl-user-fragment--rx-4
  - id: systemctl-user-fragment--rx-5
