# objectives/persistence/system/service/systemd/obfuscated-bash-builder.yaml
# Bash systemd-user persistence builder embedded in obfuscated JS/script.
# JS-string-array obfuscators split a multi-line bash script into many
# small fragments, each rotated through a decoder. The reassembled
# script writes a [Service] unit to ~/.config/systemd/user/, sets
# Restart=always + WantedBy=default.target, and `systemctl --user
# enable && start`s it. The individual fragments here are uncommon
# in legitimate code; the cluster is hostile.
#
# Examples (npm xmlrpc, xmrdropper variants):
#   `\nSYSTEMD_PATH=$HOME/.config/systemd/user`
#   `Restart=always\nWantedBy=default.target`
#   `systemctl --user enable $APPNAME.service`
defaults:
  platforms: [unix]
  for: [scripts, source]
  attack: T1543.002
  mbc: F0011
traits:
- id: systemd-path-fragment
  desc: systemd user service path fragment
  crit: suspicious
  conf: 0.92
  if:
    type: text
    regex: '\\x0aSYSTEMD|SYSTEMD_PATH=\$HOME|systemd/us|systemd[\\/]+user[\\/]+'
- id: systemctl-user-fragment--rx-1
  desc: Regex component marker
  crit: component
  conf: 0.92
  if:
    type: text
    regex: '\bsystemctl\s+--user\b'
- id: systemctl-user-fragment--rx-2
  desc: Regex component marker
  crit: component
  conf: 0.92
  if:
    type: text
    regex: 'systemctl\\\\x20--u'
- id: systemctl-user-fragment--rx-3
  desc: Regex component marker
  crit: component
  conf: 0.92
  if:
    type: text
    regex: '--user\\\\x20enab'
- id: systemctl-user-fragment--rx-4
  desc: Regex component marker
  crit: component
  conf: 0.92
  if:
    type: text
    regex: 'user\\\\x20stop'
- id: systemctl-user-fragment--rx-5
  desc: Regex component marker
  crit: component
  conf: 0.92
  if:
    type: text
    regex: 'user\\\\x20start'
- id: service-unit-write-fragment
  desc: Service unit write redirection fragment
  crit: suspicious
  conf: 0.9
  if:
    type: text
    regex: '\.service\\\\x20>|\.service\s*>|E\.service\\\\x20'
- id: restart-always-fragment
  desc: Restart=always systemd directive fragment
  crit: suspicious
  conf: 0.92
  if:
    type: text
    regex: 'Restart=always|start=alwa|rt=alwa'
- id: wanted-by-default-target-fragment
  desc: WantedBy default.target fragment
  crit: suspicious
  conf: 0.92
  if:
    type: text
    regex: 'WantedBy=default\.target|t\.target\\\\x0a|=defaul.{0,5}target'
- id: bash-os-release-read
  desc: Reads /etc/os-release for arch detection
  crit: suspicious
  conf: 0.85
  if:
    type: text
    regex: '/etc/os-re|/etc/os-release|os-release'
composite_rules:
- id: obfuscated-systemd-bash-builder
  desc: Embedded systemd persistence bash builder
  crit: hostile
  conf: 0.96
  needs: 3
  any:
  - id: systemd-path-fragment
  - id: systemctl-user-fragment
  - id: service-unit-write-fragment
  - id: restart-always-fragment
  - id: wanted-by-default-target-fragment
  - id: bash-os-release-read
- id: systemctl-user-fragment
  desc: systemctl --user command fragment
  crit: suspicious
  conf: 0.92
  any:
  - id: systemctl-user-fragment--rx-1
  - id: systemctl-user-fragment--rx-2
  - id: systemctl-user-fragment--rx-3
  - id: systemctl-user-fragment--rx-4
  - id: systemctl-user-fragment--rx-5
  needs: 1
