# objectives/persistence/system/launchd/interpreter/shell-macos.yaml
defaults:
  platforms: [unix]
  for: [shell]
  attack: T1543.004

traits:
  - id: shell-launchagents-plist-write
    desc: Shell writes LaunchAgent plist
    crit: suspicious
    conf: 0.94
    if:
      type: text
      regex: 'Library/LaunchAgents[\s\S]{0,420}<key>RunAtLoad</key>'

  - id: shell-launchdaemons-plist-write
    desc: Shell writes LaunchDaemon plist
    crit: suspicious
    conf: 0.94
    if:
      type: text
      regex: 'Library/LaunchDaemons[\s\S]{0,420}<key>RunAtLoad</key>'

  - id: shell-launchctl-load-plist
    desc: Shell loads launchd plist
    crit: suspicious
    conf: 0.94
    if:
      type: text
      regex: 'launchctl load [^;\n]{0,160}\.plist'

  - id: shell-launchd-start-interval
    desc: Shell launchd StartInterval
    crit: suspicious
    conf: 0.92
    if:
      type: text
      regex: '<key>StartInterval</key>[\s\S]{0,80}<integer>'

composite_rules:
  - id: shell-macos-launchd-persistence
    desc: Shell macOS launchd persistence
    crit: hostile
    conf: 0.98
    any:
      - id: shell-launchagents-plist-write
      - id: shell-launchdaemons-plist-write
    all:
      - id: shell-launchctl-load-plist

  - id: shell-macos-scheduled-launchd
    desc: Shell macOS scheduled launchd
    crit: hostile
    conf: 0.97
    all:
      - id: shell-launchd-start-interval
      - id: shell-launchctl-load-plist
    any:
      - id: shell-launchagents-plist-write
      - id: shell-launchdaemons-plist-write

  - id: shell-macos-launchd-var-persistence
    desc: Shell macOS launchd var persistence
    crit: hostile
    conf: 0.97
    all:
      - id: micro-behaviors/os/service/control::launchctl-load-variable
    any:
      - id: shell-launchagents-plist-write
      - id: shell-launchdaemons-plist-write
      - id: shell-launchd-start-interval

  - id: shell-launchd-vm-aware-daemon-bundle
    desc: Shell installs launchd persistence for VM-aware daemon payload
    crit: hostile
    conf: 0.97
    platforms: [macos]
    for: [archives]
    scope: archive
    all:
      - id: shell-launchctl-load-plist
      - id: objectives/anti-analysis/emulator-detect/hardware::qemu-string
    any:
      - id: objectives/anti-analysis/emulator-detect/hardware::vm-hardware
      - id: micro-behaviors/communications/socket/listen/
      - id: micro-behaviors/communications/socket/connect/
      - id: objectives/persistence/system/daemon/init::unix-daemon-persistence
