# objectives/supply-chain/hidden-payload/extensions/vscode-malicious.yaml
defaults:
  platforms:
    - unix
    - windows
  for: [scripts, source, documents]
traits:
  - id: api-usage
    desc: Uses VSCode API (extension marker)
    crit: notable
    conf: 1.0
    for:
      - javascript
      - typescript
    if:
      id: micro-behaviors/process/create/shell/host::vscode-extension-marker
  - id: postInstall-word
    desc: postInstall keyword
    crit: notable
    conf: 0.7
    for:
      - javascript
      - typescript
    if:
      type: string_literal
      substr: postInstall
  - id: exec-pattern
    desc: exec function pattern
    crit: component
    conf: 0.7
    for:
      - javascript
      - typescript
    unless:
      - type: basename
        regex: \.map$
      - type: basename
        regex: (^|[\\/])app\.asar\.unpacked[\\/]dist[\\/]
    if:
      type: text
      regex: exec(Sync|Async)\(
    not:
      - substr: node_modules/
      - substr: xo && node test.js
  - id: readFile-pattern
    desc: readFile function pattern
    crit: component
    conf: 0.6
    for:
      - javascript
      - typescript
    if:
      type: string_literal
      regex: read(File|FileSync)
  - id: vscode-namespace
    desc: Uses VSCode namespace
    crit: notable
    conf: 1.0
    for:
      - javascript
      - typescript
    if:
      id: micro-behaviors/process/create/shell/host::namespace
  - id: single-char-names-dup
    desc: High ratio of single character function names
    crit: baseline
    conf: 0.7
    for: [javascript, typescript]
    unless:
      - id: metadata/library/structure::js-webpack-bundle
      - id: metadata/library/structure::js-large-module-bundle
      - id: metadata/package/tooling::webpack-bundled
      - type: basename
        regex: \.test\.ts$
    if:
      type: metrics
      field: functions.single_char_names
      min: 10
  - id: file-copy-ops
    desc: File copy operations
    crit: notable
    conf: 1.0
    for:
      - javascript
      - typescript
    if:
      type: string_literal
      regex: copy(File|FileSync)
  - id: vscode-workspace-state
    desc: Accesses VSCode workspaceState (infection tracking)
    crit: notable
    conf: 0.8
    for: [javascript, typescript]
    if:
      type: string_literal
      substr: workspaceState
  - id: source-map-marker
    desc: Explicit source map marker in bundled output
    crit: baseline
    conf: 0.9
    for: [javascript, typescript]
    if:
      type: text
      substr: sourceMappingURL=
composite_rules:
  - id: vscode-decoy-source-map
    desc: VSCode extension with clean source and obfuscated output divergence
    crit: suspicious
    conf: 0.95
    for: [javascript, typescript]
    all:
      - id: single-char-names-dup
      - id: vscode-namespace
    unless:
      - id: metadata/library/structure::js-webpack-bundle
      - id: metadata/library/structure::js-large-module-bundle
      - id: source-map-marker
      - id: metadata/package::terser-output
      - id: metadata/package/tooling::esbuild-bundled
      - id: metadata/package/tooling::webpack-bundled
  - id: vscode-persistence-state
    desc: VSCode extension tracking state across sessions (potential persistence)
    crit: suspicious
    conf: 0.9
    for: [javascript, typescript]
    all:
      - id: vscode-workspace-state
      - id: micro-behaviors/process/create/shell/host::activate-function
  - id: vscode-remote-payload-exec
    desc: VSCode extension fetching and executing remote payload (blockchain C2)
    crit: hostile
    conf: 0.98
    mbc: B0024
    attack: T1102
    for: [javascript, typescript]
    all:
      - id: objectives/command-and-control/infrastructure/blockchain::blockchain-c2-fetch
      - id: micro-behaviors/process/interpreter/eval::js-eval-raw
  - id: hostile-geofencing
    desc: Evasive geofencing targeting restricted regions (CIS/Iran/DPRK)
    crit: hostile
    conf: 0.99
    attack: T1614
    for: [javascript, typescript]
    all:
      - id: objectives/anti-analysis/geofencing/location::restricted-region-geofencing
      - id: objectives/anti-analysis/geofencing/location::js-timezone-detection
  - id: advanced-russian-geofencing
    desc: Advanced Russian region geofencing (TZ + Locale + JS Intl)
    crit: hostile
    conf: 0.99
    attack: T1614
    for: [javascript, typescript]
    all:
      - id: objectives/anti-analysis/geofencing/location::cis-region-indicator
      - id: objectives/anti-analysis/geofencing/location::js-timezone-detection
      - id: objectives/anti-analysis/geofencing/location::func-name-russian-specific
  - id: state-markers
    desc: State management markers (globalState or postInstall)
    crit: notable
    conf: 0.75
    for:
      - javascript
      - typescript
    any:
      - id: micro-behaviors/process/create/shell/host::global-state
      - id: postInstall-word
  - id: credential-paths
    desc: Sensitive credential path references
    crit: notable
    conf: 0.8
    for:
      - javascript
      - typescript
    any:
      - id: micro-behaviors/fs/path/sensitive/ssh::keys
      - id: objectives/credential-access/cloud/token::aws-config-dup
      - id: objectives/credential-access/env/secrets::dot-npmrc
      - id: metadata/package/config::dotenv-file
  - id: send-methods
    desc: HTTP send methods
    crit: notable
    conf: 0.7
    for:
      - javascript
      - typescript
    any:
      - id: micro-behaviors/communications/http/post/
      - id: micro-behaviors/communications/http/request/
      - id: micro-behaviors/communications/http/lib/axios::axios
  - id: oauth-uri-exfil
    desc: VSCode URI OAuth code exfiltration
    crit: hostile
    conf: 0.98
    mbc: E1056
    attack: T1041
    for:
      - javascript
      - typescript
    all:
      - id: objectives/supply-chain/metadata-anomaly/manifest/vscode::oauth-code-capture
      - id: micro-behaviors/communications/http/request::https-request
      - id: micro-behaviors/communications/http/request::post-method-string
      - id: objectives/exfiltration/oob::webhook-site
  - id: post-install-exec
    desc: Post-install script execution (dropper pattern)
    crit: suspicious
    conf: 0.95
    mbc: B0024
    attack: T1059
    for:
      - javascript
      - typescript
    all:
      - id: objectives/supply-chain/metadata-anomaly/manifest/vscode::activate-comp
      - id: state-markers
      - id: exec-pattern
  - id: auto-activated-reverse-shell
    desc: Auto-activated VSCode reverse shell
    crit: hostile
    conf: 0.995
    mbc: B0022
    attack: T1059
    for:
      - javascript
      - typescript
    all:
      - id: micro-behaviors/process/create/shell/host::vscode-extension-marker
      - id: objectives/supply-chain/metadata-anomaly/manifest/vscode::activate-comp
      - id: objectives/command-and-control/reverse-shell/socket::cmd-shell-rfc1918-bridge
  - id: credential-theft-comp
    desc: VSCode extension accessing sensitive credentials
    crit: suspicious
    conf: 0.9
    mbc: E1519
    attack: T1555
    for:
      - javascript
      - typescript
    all:
      - id: credential-paths
      - id: readFile-pattern
      - id: send-methods
    downgrade:
      any:
        - id: objectives/anti-static/obfuscation/code-metrics/structure::minified-code
        - id: metadata/package/tooling::esbuild-bundled
  - id: persistence-locations
    desc: References persistence locations (startup/autorun paths)
    crit: notable
    conf: 0.8
    for:
      - javascript
      - typescript
    any:
      - id: micro-behaviors/os/autorun::startup-keyword
      - id: micro-behaviors/os/autorun/scheduled::crontab-keyword
      - id: micro-behaviors/os/autorun::launchd-keyword
      - id: objectives/persistence/system/service/systemd::systemd-keyword
  - id: persistence
    desc: Extension establishing persistence
    crit: suspicious
    conf: 0.95
    mbc: F0011
    attack: T1547
    for:
      - javascript
      - typescript
    all:
      - id: api-usage
      - id: persistence-locations
      - id: micro-behaviors/fs/write/file::js-file-write
  - id: file-write-ops
    desc: File write or copy operations
    crit: notable
    conf: 1.0
    for:
      - javascript
      - typescript
    any:
      - id: micro-behaviors/fs/write/file::js-file-write
      - id: file-copy-ops
