# objectives/supply-chain/hidden-payload/extensions/vscode/vscode-malicious.yaml
defaults:
  platforms:
  - unix
  - windows
  for: [scripts, source, documents]
traits:
- id: postInstall-word
  desc: postInstall keyword
  crit: notable
  conf: 0.7
  for:
  - javascript
  - typescript
  if:
    type: string_literal
    substr: postInstall
- id: exec-pattern
  desc: exec function pattern
  crit: component
  conf: 0.7
  for:
  - javascript
  - typescript
  unless:
  - type: basename
    regex: (^|[\\/])app\.asar\.unpacked[\\/]dist[\\/]
  if:
    type: text
    regex: exec(Sync|Async)\(
  not:
  - substr: node_modules/
  - substr: xo && node test.js
# single-char function-name ratio is a generic minification/obfuscation metric
# — it lives canonically in objectives/anti-static/obfuscation/code-metrics/
# functions::single-char-function-names. The vscode-decoy-source-map composite
# below references it cross-tier rather than duplicating it into this
# vscode-specific directory (where it mislabeled any minified JS/TS).
- id: vscode-workspace-state
  desc: Accesses VSCode workspaceState (infection tracking)
  crit: notable
  conf: 0.8
  for: [javascript, typescript]
  if:
    type: string_literal
    substr: workspaceState
composite_rules:
- id: vscode-decoy-source-map
  desc: VSCode extension with clean source and obfuscated output divergence
  crit: suspicious
  conf: 0.95
  for: [javascript, typescript]
  all:
  - id: objectives/anti-static/obfuscation/code-metrics/functions::single-char-function-names
  - id: micro-behaviors/process/create/shell/host::namespace
  unless:
  - id: metadata/library/structure::js-webpack-bundle
  - id: metadata/library/structure::js-large-module-bundle
  - id: metadata/package/tooling::has-sourcemap
  - id: metadata/package/quality::terser-output
  - id: metadata/package/tooling::esbuild-bundled
  - id: metadata/package/tooling::webpack-bundled
- id: vscode-persistence-state
  desc: VSCode extension tracking state across sessions (potential persistence)
  crit: suspicious
  conf: 0.9
  for: [javascript, typescript]
  all:
  - id: vscode-workspace-state
  - id: micro-behaviors/process/create/shell/host::activate-function
- id: vscode-remote-payload-exec
  desc: VSCode extension fetching and executing remote payload (blockchain C2)
  crit: hostile
  conf: 0.98
  mbc: B0024
  attack: T1102
  for: [javascript, typescript]
  all:
  - id: objectives/command-and-control/infrastructure/blockchain/multi-chain::blockchain-c2-fetch
  - id: micro-behaviors/process/interpreter/eval::js-eval-raw
- id: vscode-http-response-eval
  desc: VSCode extension evals HTTP response
  crit: hostile
  conf: 0.99
  attack: T1059.007
  mbc: B0024
  for:
  - javascript
  - typescript
  near_bytes: 4096
  all:
  - id: micro-behaviors/process/create/shell/host::vscode-extension-marker
  - id: objectives/command-and-control/dropper/delivery/fetch-eval::response-end-then-eval
  any:
  - id: micro-behaviors/communications/http/request/client::http-get-dup
  - id: metadata/import/builtin::require-http
- id: hostile-geofencing
  desc: Evasive geofencing targeting restricted regions (CIS/Iran/DPRK)
  crit: hostile
  conf: 0.99
  attack: T1614
  for: [javascript, typescript]
  all:
  - id: objectives/anti-analysis/geofencing/region/cluster::restricted-region-geofencing
  - id: objectives/anti-analysis/geofencing/timezone::js-timezone-detection
- id: advanced-russian-geofencing
  desc: Advanced Russian region geofencing (TZ + Locale + JS Intl)
  crit: hostile
  conf: 0.99
  attack: T1614
  for: [javascript, typescript]
  all:
  - id: objectives/anti-analysis/geofencing/region/cluster::cis-region-indicator
  - id: objectives/anti-analysis/geofencing/timezone::js-timezone-detection
  - id: objectives/anti-analysis/geofencing/region/cluster::func-name-russian-specific
- id: state-markers
  desc: State management markers (globalState or postInstall)
  crit: notable
  conf: 0.75
  for:
  - javascript
  - typescript
  any:
  - id: micro-behaviors/process/create/shell/host::global-state
  - id: postInstall-word
- id: credential-paths
  desc: Sensitive credential path references
  crit: notable
  conf: 0.8
  for:
  - javascript
  - typescript
  any:
  - id: micro-behaviors/fs/path/sensitive/ssh::keys
  - id: objectives/credential-access/cloud/token::aws-config-dup
  - id: objectives/credential-access/env/secrets/named-access::dot-npmrc
  - id: metadata/package/config::dotenv-file
- id: send-methods
  desc: HTTP send methods
  crit: notable
  conf: 0.7
  for:
  - javascript
  - typescript
  any:
  - id: micro-behaviors/communications/http/post/
  - id: micro-behaviors/communications/http/request/
  - id: micro-behaviors/communications/http/lib/axios::axios
- id: oauth-uri-exfil
  desc: VSCode URI OAuth code exfiltration
  crit: hostile
  conf: 0.98
  mbc: E1056
  attack: T1041
  for:
  - javascript
  - typescript
  all:
  - id: objectives/supply-chain/metadata-anomaly/manifest/vscode::oauth-code-capture
  - id: micro-behaviors/communications/http/request/client::https-request
  - id: micro-behaviors/communications/http/request/client::post-method-string
  - id: objectives/exfiltration/oob/endpoint::webhook-site
- id: post-install-exec
  desc: Post-install script execution (dropper pattern)
  crit: suspicious
  conf: 0.95
  mbc: B0024
  attack: T1059
  for:
  - javascript
  - typescript
  all:
  - id: objectives/supply-chain/metadata-anomaly/manifest/vscode::activate-comp
  - id: state-markers
  - id: exec-pattern
- id: auto-activated-reverse-shell
  desc: Auto-activated VSCode reverse shell
  crit: hostile
  conf: 0.995
  mbc: B0022
  attack: T1059
  for:
  - javascript
  - typescript
  all:
  - id: micro-behaviors/process/create/shell/host::vscode-extension-marker
  - id: objectives/supply-chain/metadata-anomaly/manifest/vscode::activate-comp
  - id: objectives/command-and-control/reverse-shell/dup::cmd-shell-rfc1918-bridge
- id: credential-theft-comp
  desc: VSCode extension accessing sensitive credentials
  crit: suspicious
  conf: 0.9
  mbc: E1519
  attack: T1555
  for:
  - javascript
  - typescript
  all:
  - id: credential-paths
  - id: micro-behaviors/fs/file/read::read-node
  - id: send-methods
  unless:
  - id: metadata/package/testing/presence/fixture::test-filename-pattern
  - id: metadata/package/testing/harness/runtime::is-test
  downgrade:
    any:
    - id: objectives/anti-static/obfuscation/code-metrics/structure::minified-code
    - id: metadata/package/tooling::esbuild-bundled
- id: persistence-locations
  desc: References persistence locations (startup/autorun paths)
  crit: suspicious
  conf: 0.8
  for:
  - javascript
  - typescript
  any:
  - id: micro-behaviors/os/autorun/keyword
  - id: micro-behaviors/os/autorun/scheduled::crontab-keyword
  - id: objectives/persistence/system/service/systemd::systemd-keyword
- id: persistence
  desc: Extension establishing persistence
  crit: suspicious
  conf: 0.95
  mbc: F0011
  attack: T1547
  for:
  - javascript
  - typescript
  all:
  - id: micro-behaviors/process/create/shell/host::vscode-extension-marker
  - id: persistence-locations
  - id: micro-behaviors/fs/write/file/direct::js-file-write
- id: file-write-ops
  desc: File write or copy operations
  crit: notable
  conf: 1.0
  for:
  - javascript
  - typescript
  any:
  - id: micro-behaviors/fs/write/file/direct::js-file-write
  - id: micro-behaviors/fs/file/copy::fs-copy
